# etc/hotplug.d/dsl/10_trigger_xtm
if [ "$ACTION" = "ifup" ] ; then
	/etc/init.d/xtm start
	[ -x "/etc/init.d/vopi_intfs" ] && /etc/init.d/vopi_intfs restart
else
	[ -x "/etc/init.d/vopi_intfs" ] && /etc/init.d/vopi_intfs stop
	/etc/init.d/xtm stop
	[ -x "/etc/init.d/vopi_intfs" ] && /etc/init.d/vopi_intfs start
fi
