# etc/hotplug.d/dsl/10_trigger_xtm
[ ifup = "$ACTION" ] || exit
/etc/init.d/xtm enabled || exit 0
/etc/init.d/xtm start
