# etc/hotplug.d/dsl/10_trigger_xtm
if [ "$ACTION" = "ifup" ] ; then
        /etc/init.d/xtm start
else
        /etc/init.d/xtm stop
fi
