# etc/hotplug.d/platform/00-sfp
. /lib/functions.sh

name="sfp"

if [ "$NAME" = "$name" ]; then
	if [ "$ACTION" = "plugin" ]; then
		ubus send "sfp" '{"status":"plugin"}'
	elif [ "$ACTION" = "unplug" ]; then
		ubus send "sfp" '{"status":"unplug"}'
	elif [ "$ACTION" = "linkup" ]; then
		ubus send "sfp" '{"status":"linkup"}'
	elif [ "$ACTION" = "linkdown" ]; then
		ubus send "sfp" '{"status":"linkdown"}'
	else
		echo "invalid action message get from Uevent"
	fi
fi
