# etc/hotplug.d/iface/40-mobiled
ubus send mobiled.network "{\"interface\":\"$INTERFACE\",\"action\":\"$ACTION\"}"
