# etc/hotplug.d/iface/99-kernelpanic
#!/bin/sh

[ -e "/proc/prozone/panic" ] || exit 0
. "$IPKG_INSTROOT/lib/functions.sh"
[ "$ACTION" = "ifup" ] && [ "$DEVICE" = "$(uci_get network wan ifname)" ] && sh /sbin/kernel-panic-handler
