# etc/hotplug.d/iface/99-warmboot
[ ifup = "$ACTION" ] || exit

. $IPKG_INSTROOT/lib/config/uci.sh

[ "$(uci_get_state system warmboot sent 0)" != "1" ] || exit

. $IPKG_INSTROOT/lib/functions.sh
. $IPKG_INSTROOT/lib/functions/syslog.sh
. $IPKG_INSTROOT/lib/functions/reboot_reason.sh
Haswan=$(uci get network.wan)
if [ -z $Haswan ];then
        syslog_iface=""

        get_syslog_iface syslog_iface lan

        [ "$syslog_iface" == "$INTERFACE" ] && {
                reboot=$(uci_get_state system warmboot reboot)
                reboot_reason_str=$(get_reboot_reason $reboot)
                logger -t warmboot "reason:${reboot} ${reboot_reason_str}"
                uci_set_state system warmboot sent 1
        }
else
        syslog_iface=""

        get_syslog_iface syslog_iface

        [ "$syslog_iface" == "$INTERFACE" ] && {
                reboot=$(uci_get_state system warmboot reboot)
                reboot_reason_str=$(get_reboot_reason $reboot)
                logger -t warmboot "reason:${reboot} ${reboot_reason_str}"
                uci_set_state system warmboot sent 1
        }

fi

