# chroot/etc/boards/VCNT-H/config/etc/uci-defaults/tch_0010-remove-usb3
#!/bin/sh
#remove file for boards that do not support usb3
[ -f /lib/modules/$(uname -r)/xhci-hcd.ko ] && rm -rf /lib/modules/$(uname -r)/xhci-hcd.ko
