# etc/hotplug.d/ntp/05-hostmanager
#!/bin/sh

[ "$ACTION" = stratum ] || [ "$ACTION" = step ] || exit 0

#inform hostmanager
ubus send ntp.connected
