# etc/hotplug.d/ntp/20-bulkdata
#!/bin/sh

[ "$ACTION" = stratum ] || [ "$ACTION" = step ] || exit 0
#After time is synced, restart bulkdata to update /etc/crontab/root
/etc/init.d/bulkdata restart
