# etc/boards/VANT-W/config/etc/uci-defaults/tch_0005-slic
#!/bin/sh

#for now, nothing to be done for other platforms ...
