# etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

##### install new ADB agent and remove ipk
/bin/opkg install --force-reinstall /usr/share/homeware-ipks/adb-speedtest_2.0_arm_armv7-a.ipk && rm /usr/share/homeware-ipks/adb-speedtest_2.0_arm_armv7-a.ipk

##### avoid further execution
chmod -x /etc/rc.local


exit 0
