# etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

. /lib/functions.sh
. /lib/functions/TIM_lib.sh

config_load "env"
enabled
config_get_bool enabled var unlockedstatus 0

if [ $enabled -eq 0 ]
then
	TIM_mount
fi

. /lib/functions/run_main.sh

exit 0
