# etc/uci-defaults/70-464xlat_skiplog
#!/bin/sh

echo "ip6tables -t mangle -I FORWARD -o 464-wan6_4 -j SKIPLOG" >> /etc/firewall.user
echo "ip6tables -t mangle -I FORWARD -i 464-wan6_4 -j SKIPLOG" >> /etc/firewall.user

exit 0
