# etc/uci-defaults/70-mapt_skiplog
#!/bin/sh

echo "ip6tables -t mangle -I FORWARD -o map-mapt -j SKIPLOG" >> /etc/firewall.user
echo "ip6tables -t mangle -I FORWARD -i map-mapt -j SKIPLOG" >> /etc/firewall.user

exit 0
