# etc/uci-defaults/conf-web-add-parental
--uci -q del_list web.ruleset_main.rules=parentalmodal
--uci add_list web.ruleset_main.rules=parentalmodal

--uci -q delete web.parentalmodal
--uci set web.parentalmodal=rule
--uci set web.parentalmodal.target=/modals/parental-modal.lp
--uci add_list web.parentalmodal.roles=admin
--uci add_list web.parentalmodal.roles=engineer

uci -q delete web.ruleset_intercept
uci set web.ruleset_intercept=ruleset
uci -q del_list web.ruleset_intercept.rules=parentalblock
uci add_list web.ruleset_intercept.rules=parentalblock

uci -q delete web.parentalblock
uci set web.parentalblock=rule
uci set web.parentalblock.target=/parental-block.lp
uci add_list web.parentalblock.roles=intercept

uci -q delete web.intercept
uci set web.intercept=sessionmgr
uci set web.intercept.cookiepath=/
uci set web.intercept.authpath=/authenticate
uci set web.intercept.passpath=/password
uci set web.intercept.loginpath=/login.lp
uci set web.intercept.timeout=30
uci set web.intercept.ruleset=ruleset_intercept
uci set web.intercept.default_user=usr_intercept
uci add_list web.intercept.users=usr_intercept

uci -q delete web.usr_intercept
uci set web.usr_intercept=user
uci set web.usr_intercept.name='intercept'
uci set web.usr_intercept.role='intercept'
uci commit web
