# etc/uci-defaults/conf-web-add-qos
uci -q del_list web.ruleset_main.rules=qosqueuemodal
uci add_list web.ruleset_main.rules=qosqueuemodal

uci -q delete web.qosqueuemodal
uci set web.qosqueuemodal=rule
uci set web.qosqueuemodal.target=/modals/qos-queue-modal.lp
uci add_list web.qosqueuemodal.roles=admin

uci -q del_list web.ruleset_main.rules=qosclassifymodal
uci add_list web.ruleset_main.rules=qosclassifymodal

uci -q delete web.qosclassifymodal
uci set web.qosclassifymodal=rule
uci set web.qosclassifymodal.target=/modals/qos-classify-modal.lp
uci add_list web.qosclassifymodal.roles=admin

uci commit web
