# etc/uci-defaults/conf-web-add-mobiledongle
uci -q del_list web.ruleset_main.rules=mobiledongleinformationmodal
uci add_list web.ruleset_main.rules=mobiledongleinformationmodal

uci -q delete web.mobiledongleinformationmodal
uci set web.mobiledongleinformationmodal=rule
uci set web.mobiledongleinformationmodal.target=/modals/mobiledongle-information-modal.lp
uci add_list web.mobiledongleinformationmodal.roles=admin
uci add_list web.mobiledongleinformationmodal.roles=engineer

uci -q del_list web.ruleset_main.rules=mobiledongleconfigurationmodal
uci add_list web.ruleset_main.rules=mobiledongleconfigurationmodal

uci -q delete web.mobiledongleconfigurationmodal
uci set web.mobiledongleconfigurationmodal=rule
uci set web.mobiledongleconfigurationmodal.target=/modals/mobiledongle-configuration-modal.lp
uci add_list web.mobiledongleconfigurationmodal.roles=admin
uci add_list web.mobiledongleconfigurationmodal.roles=engineer
