# etc/uci-defaults/conf-web-add-contentsharing
uci -q del_list web.ruleset_main.rules=contentsharing
uci add_list web.ruleset_main.rules=contentsharing

uci -q delete web.contentsharing
uci set web.contentsharing=rule
uci set web.contentsharing.target=/modals/contentsharing-modal.lp
uci add_list web.contentsharing.roles=admin
uci commit web
