# etc/uci-defaults/conf-web-add-ethwan
uci -q del_list web.ruleset_main.rules=ethwanmodal
uci add_list web.ruleset_main.rules=ethwanmodal

uci -q delete web.ethwanmodal
uci set web.ethwanmodal=rule
uci set web.ethwanmodal.target=/modals/ethwan-modal.lp
uci add_list web.ethwanmodal.roles=guest
uci add_list web.ethwanmodal.roles=admin
uci add_list web.ethwanmodal.roles=engineer

