# etc/uci-defaults/conf-web-add-fon
uci -q del_list web.ruleset_main.rules=fon
uci add_list web.ruleset_main.rules=fon

uci -q delete web.fon
uci set web.fon=rule
uci set web.fon.target=/modals/fon-modal.lp
uci add_list web.fon.roles=admin
uci commit web
