# etc/uci-defaults/conf-web-add-nfc
uci -q del_list web.ruleset_main.rules=nfc
uci add_list web.ruleset_main.rules=nfc

uci -q delete web.nfc
uci set web.nfc=rule
uci set web.nfc.target=/modals/nfc-modal.lp
uci add_list web.nfc.roles=admin
