# etc/uci-defaults/conf-web-add-qeo
uci -q del_list web.ruleset_main.rules=qeoregstatus
uci add_list web.ruleset_main.rules=qeoregstatus

uci -q delete web.qeoregstatus
uci set web.qeoregstatus=rule
uci set web.qeoregstatus.target=/qeo/regstatus
uci add_list web.qeoregstatus.roles=guest
uci add_list web.qeoregstatus.roles=admin
uci add_list web.qeoregstatus.roles=engineer
