# etc/uci-defaults/conf-web-add-voice
uci -q del_list web.ruleset_main.rules=mmpbxinfomodal
uci add_list web.ruleset_main.rules=mmpbxinfomodal

uci -q delete web.mmpbxinfomodal
uci set web.mmpbxinfomodal=rule
uci set web.mmpbxinfomodal.target=/modals/mmpbx-info-modal.lp
uci add_list web.mmpbxinfomodal.roles=admin

uci -q del_list web.ruleset_main.rules=mmpbxglobalmodal
uci add_list web.ruleset_main.rules=mmpbxglobalmodal

uci -q delete web.mmpbxglobalmodal
uci set web.mmpbxglobalmodal=rule
uci set web.mmpbxglobalmodal.target=/modals/mmpbx-global-modal.lp
uci add_list web.mmpbxglobalmodal.roles=admin

uci -q del_list web.ruleset_main.rules=mmpbxprofilemodal
uci add_list web.ruleset_main.rules=mmpbxprofilemodal

uci -q delete web.mmpbxprofilemodal
uci set web.mmpbxprofilemodal=rule
uci set web.mmpbxprofilemodal.target=/modals/mmpbx-profile-modal.lp
uci add_list web.mmpbxprofilemodal.roles=admin

uci -q del_list web.ruleset_main.rules=mmpbxservicemodal
uci add_list web.ruleset_main.rules=mmpbxservicemodal

uci -q delete web.mmpbxservicemodal
uci set web.mmpbxservicemodal=rule
uci set web.mmpbxservicemodal.target=/modals/mmpbx-service-modal.lp
uci add_list web.mmpbxservicemodal.roles=admin

uci -q del_list web.ruleset_main.rules=mmpbxlogmodal
uci add_list web.ruleset_main.rules=mmpbxlogmodal

uci -q delete web.mmpbxlogmodal
uci set web.mmpbxlogmodal=rule
uci set web.mmpbxlogmodal.target=/modals/mmpbx-log-modal.lp
uci add_list web.mmpbxlogmodal.roles=admin

uci -q del_list web.ruleset_main.rules=mmpbxdectmodal
uci add_list web.ruleset_main.rules=mmpbxdectmodal

uci -q delete web.mmpbxdectmodal
uci set web.mmpbxdectmodal=rule
uci set web.mmpbxdectmodal.target=/modals/mmpbx-dect-modal.lp
uci add_list web.mmpbxdectmodal.roles=admin

uci -q del_list web.ruleset_main.rules=mmpbxcontactsmodal
uci add_list web.ruleset_main.rules=mmpbxcontactsmodal

uci -q delete web.mmpbxcontactsmodal
uci set web.mmpbxcontactsmodal=rule
uci set web.mmpbxcontactsmodal.target=/modals/mmpbx-contacts-modal.lp
uci add_list web.mmpbxcontactsmodal.roles=admin

uci -q del_list web.ruleset_main.rules=mmpbxsipdevicemodal
uci add_list web.ruleset_main.rules=mmpbxsipdevicemodal

uci -q delete web.mmpbxsipdevicemodal
uci set web.mmpbxsipdevicemodal=rule
uci set web.mmpbxsipdevicemodal.target=/modals/mmpbx-sipdevice-modal.lp
uci add_list web.mmpbxsipdevicemodal.roles=admin

uci commit web
