# etc/uci-defaults/conf-web-add-voice
#NG-78959

uci -q del_list web.ruleset_main.rules=mmpbxglobalmodal
uci add_list web.ruleset_main.rules=mmpbxglobalmodal

uci -q delete web.mmpbxglobalmodal
uci set web.mmpbxglobalmodal=rule
uci set web.mmpbxglobalmodal.target=/modals/mmpbx-global-modal.lp
#NG-78959 uci add_list web.mmpbxglobalmodal.roles=admin
uci add_list web.mmpbxglobalmodal.roles=engineer

uci -q del_list web.ruleset_main.rules=mmpbxprofilemodal
uci add_list web.ruleset_main.rules=mmpbxprofilemodal

uci -q delete web.mmpbxprofilemodal
uci set web.mmpbxprofilemodal=rule
uci set web.mmpbxprofilemodal.target=/modals/mmpbx-profile-modal.lp
#NG-78959 uci add_list web.mmpbxprofilemodal.roles=admin
uci add_list web.mmpbxprofilemodal.roles=engineer

uci -q del_list web.ruleset_main.rules=mmpbxincomingmapmodal
uci -q del_list web.ruleset_main.rules=mmpbxoutgoingmapmodal
uci -q delete web.mmpbxincomingmapmodal
uci -q delete web.mmpbxoutgoingmapmodal

uci -q del_list web.ruleset_main.rules=mmpbxinoutgoingmodal
uci add_list web.ruleset_main.rules=mmpbxinoutgoingmodal

uci -q delete web.mmpbxinoutgoingmodal
uci set web.mmpbxinoutgoingmodal=rule
uci set web.mmpbxinoutgoingmodal.target=/modals/mmpbx-inoutgoingmap-modal.lp
uci add_list web.mmpbxinoutgoingmodal.roles=admin
uci add_list web.mmpbxinoutgoingmodal.roles=engineer


uci -q del_list web.ruleset_main.rules=mmpbxservicemodal
uci add_list web.ruleset_main.rules=mmpbxservicemodal

uci -q delete web.mmpbxservicemodal
uci set web.mmpbxservicemodal=rule
uci set web.mmpbxservicemodal.target=/modals/mmpbx-service-modal.lp
#NG-78959 uci add_list web.mmpbxservicemodal.roles=admin
uci add_list web.mmpbxservicemodal.roles=engineer

uci -q del_list web.ruleset_main.rules=mmpbxlogmodal
uci add_list web.ruleset_main.rules=mmpbxlogmodal

uci -q delete web.mmpbxlogmodal
uci set web.mmpbxlogmodal=rule
uci set web.mmpbxlogmodal.target=/modals/mmpbx-log-modal.lp
uci add_list web.mmpbxlogmodal.roles=admin
uci add_list web.mmpbxlogmodal.roles=engineer

uci -q del_list web.ruleset_main.rules=mmpbxdectmodal
uci add_list web.ruleset_main.rules=mmpbxdectmodal

uci -q delete web.mmpbxdectmodal
uci set web.mmpbxdectmodal=rule
uci set web.mmpbxdectmodal.target=/modals/mmpbx-dect-modal.lp
uci add_list web.mmpbxdectmodal.roles=admin
uci add_list web.mmpbxdectmodal.roles=engineer

uci -q del_list web.ruleset_main.rules=mmpbxcontactsmodal
#uci add_list web.ruleset_main.rules=mmpbxcontactsmodal

uci -q delete web.mmpbxcontactsmodal
uci set web.mmpbxcontactsmodal=rule
uci set web.mmpbxcontactsmodal.target=/modals/mmpbx-contacts-modal.lp
#uci add_list web.mmpbxcontactsmodal.roles=admin
#uci add_list web.mmpbxcontactsmodal.roles=engineer

uci -q del_list web.ruleset_main.rules=mmpbxsipdevmodal
uci add_list web.ruleset_main.rules=mmpbxsipdevmodal

uci -q delete web.mmpbxsipdevmodal
uci set web.mmpbxsipdevmodal=rule
uci set web.mmpbxsipdevmodal.target=/modals/mmpbx-sipdevice-modal.lp
uci add_list web.mmpbxsipdevmodal.roles=admin
uci add_list web.mmpbxsipdevmodal.roles=engineer

uci -q del_list web.ruleset_main.rules=mmpbxcodecTImodal
uci add_list web.ruleset_main.rules=mmpbxcodecTImodal

uci -q delete web.mmpbxcodecTImodal
uci set web.mmpbxcodecTImodal=rule
uci set web.mmpbxcodecTImodal.target=/modals/mmpbx-codec-modal.lp
uci add_list web.mmpbxcodecTImodal.roles=admin
uci add_list web.mmpbxcodecTImodal.roles=engineer

uci -q del_list web.ruleset_main.rules=mmpbxajaxvoicetab
uci add_list web.ruleset_main.rules=mmpbxajaxvoicetab

uci -q delete web.mmpbxajaxvoicetab
uci set web.mmpbxajaxvoicetab=rule
uci set web.mmpbxajaxvoicetab.target=/ajax/voicetab.lua
uci add_list web.mmpbxajaxvoicetab.roles=admin
uci add_list web.mmpbxajaxvoicetab.roles=engineer
uci add_list web.mmpbxajaxvoicetab.roles=guest

uci commit web