# etc/uci-defaults/conf-web-add-wizard
uci -q del_list web.ruleset_main.rules=wizardmodal
uci add_list web.ruleset_main.rules=wizardmodal

uci -q delete web.wizardmodal
uci set web.wizardmodal=rule
uci set web.wizardmodal.target=/modals/wizard-modal.lp
uci add_list web.wizardmodal.roles=admin

uci commit web
