# etc/uci-defaults/conf-web-add-printersharing
uci -q del_list web.ruleset_main.rules=printersharing
uci add_list web.ruleset_main.rules=printersharing

uci -q delete web.printersharing
uci set web.printersharing=rule
uci set web.printersharing.target=/modals/printersharing-modal.lp
#uci add_list web.printersharing.roles=guest
uci add_list web.printersharing.roles=admin
#uci add_list web.printersharing.roles=engineer
