# etc/uci-defaults/tch-0082-nginx-dumaos
#!/bin/sh

lua /usr/lib/netduma/netduma_rules.lua
