# etc/uci-defaults/tch-0082-nginx-dumaos
#!/bin/sh

NGINX_CONF=/etc/nginx/nginx.conf

sed -i '/ui_server/a\        include dumaos.conf;' $NGINX_CONF
