# etc/uci-defaults/tch-0082-procmanager-dumaos
#!/bin/sh
uci set dumaos.tr69=tr69; uci set dumaos.tr69.dumaos_enabled=0;uci commit dumaos
exit 0
