# etc/uci-defaults/tch_0031-pppoe-relay
#!/bin/sh

uci add_list network.lan.pppoerelay="$(uci get network.wan.ifname)"

uci commit network

exit 0
