# etc/uci-defaults/tch_8901_iperf
#!/bin/sh
#GHG-3935

uci set iperf.iperf.proto='tcp'
uci set iperf.iperf.dest_port='5001'
uci set iperf.iperf.enable='1'

uci commit iperf
