# etc/uci-defaults/tch_0090-system
#!/bin/sh

# Disable console logging
uci set system.@system[0].conloglevel=1
#uci set system.@system[0].hostname=$(echo -n $(uci get env.var.prod_friendly_name)_$(uci get env.var.local_wifi_mac |grep -o '..$'))
uci set system.@system[0].hostname=$(echo -n $(uci get env.var.prod_friendly_name))
uci set system.@system[0].log_prefix="$(uci get env.rip.factory_id)$(uci get env.rip.serial)"
uci commit system
