# etc/uci-defaults/tch_0020-network-lan
#!/bin/sh

. /lib/functions/uci-defaults.sh

ucidef_set_interface_lan "eth0_0 eth1_0 eth2_0 eth3_0"

uci set network.lan.ipaddr='192.168.1.254'
uci set network.lan.ip6assign=64
uci set network.lan.force_link=0
uci set network.lan.ip6class=6rd
uci set network.lan.ipv6=1
uci set network.lan.ip6hint=0
uci set network.lan.ip6weight=1
uci set network.lan.stp=1
uci set network.lan.rpfilter=1

uci set network.Guest=interface
uci set network.Guest.proto=static
uci set network.Guest.netmask=255.255.255.0
uci set network.Guest.ipaddr=192.168.78.1
uci set network.Guest.ifname=wl0_2
uci set network.Guest.ip6assign=64
uci set network.Guest.ip6class=6rd
uci set network.Guest.ipv6=1
uci set network.Guest.ip6hint=1
uci set network.Guest.ip6weight=0
uci set network.Guest.force_link=0
uci set network.Guest.rpfilter=1
uci add_list network.Guest.pppoerelay="$(uci get network.wan.ifname)"

ucidef_set_switch "bcmsw" "1" "0"
uci set network.bcmsw.qosimppauseenable=0
uci set network.bcmsw.jumbo=0
uci commit network

exit 0
