# etc/uci-defaults/tch_0071-ssh-pubkey
#!/bin/sh

. /lib/functions/provision.sh

# create user with root UID
create_user tchroot tchroot 0

# lock password login to allow only pub key login
passwd -l tchroot


