# etc/uci-defaults/tch_0030-network-wan
#!/bin/sh

. /lib/functions/uci-defaults.sh

uci set network.wan=interface
uci set network.wan.proto='dhcp'
uci set network.wan.reqopts='1 3 6 15 33 42 43 51 121 249 212'
uci set network.wan.vendorid=`uci get env.var.dhcp_option60`
uci set network.wan.ifname='veip0_1'
uci set network.wan.authfail='0'
uci set network.wan.iface6rd='6rd'
uci set network.wan.zone6rd='wan'
uci set network.wan.mtu='1500'
uci set network.wan.initboot='1'
uci set network.wan.ipv6='0'
uci set network.wan.auto='1'
uci set network.wan.notused=`uci get env.rip.wifi_mac`

uci set network.wwan=interface
uci set network.wwan.auto='1'

uci commit network

uci show dmordering >/dev/null 2>/dev/null || touch /etc/config/dmordering

uci set dmordering.IPInterface=ordering
uci add_list dmordering.IPInterface.order=loopback
uci add_list dmordering.IPInterface.order=lan
uci add_list dmordering.IPInterface.order=Guest
uci add_list dmordering.IPInterface.order=wan
uci add_list dmordering.IPInterface.order=wwan

uci set dmordering.dhcpv4=ordering
uci set dmordering.dhcpv4.name=Device.DHCPv4.Server.Pool.{i}.
uci add_list dmordering.dhcpv4.order=lan
uci add_list dmordering.dhcpv4.order=Guest

uci commit dmordering

exit 0
