# etc/uci-defaults/tch_9000-parameter_conversion
#!/bin/sh

exec /usr/lib/parameter_conversion/parameter_conversion.sh 2>&1
