# etc/uci-defaults/tch_9005-workaround
cp /rom/usr/lib/lua/ledframework/ubus.lua /usr/lib/lua/ledframework/ubus.lua
