# etc/uci-defaults/tch_0099-wifi_conductor
#!/bin/sh

uci set wifi_conductor.debug.cloud_home_id=`uci get env.var.serial`
uci commit wifi_conductor
