# etc/uci-defaults/tch_0020-network-lan
#!/bin/sh

. /lib/functions/uci-defaults.sh

ucidef_set_interface_lan "eth0 eth1 eth2 eth3"

#guest wifi 1
uci set network.guest1=interface
uci set network.guest1.ifname=wl0_1
uci set network.guest1.proto=static
uci set network.guest1.ipaddr=192.168.25.1
uci set network.guest1.netmask=255.255.255.0
uci set network.guest1.ip6assign=64
uci set network.guest1.ip6hint=1
uci set network.guest1.force_link=0

#guest wifi 2
uci set network.guest2=interface
uci set network.guest2.ifname=wl1_1
uci set network.guest2.proto=static
uci set network.guest2.ipaddr=192.168.26.1
uci set network.guest2.netmask=255.255.255.0
uci set network.guest2.ip6assign=64
uci set network.guest2.ip6hint=2
uci set network.guest2.force_link=0

uci set network.lan.ip6assign=64
uci set network.lan.ip6hint=0
uci set network.lan.ipv6=1
uci set network.lan.force_link=0

ucidef_set_switch "bcmsw" "1" "0"
uci set network.bcmsw.qosimppauseenable=0
uci set network.bcmsw.jumbo=0
uci commit network

exit 0
