# etc/uci-defaults/tch_0080-samba
#!/bin/sh

uci set samba.samba.name='Modem TIM'
uci set samba.samba.description='Modem TIM'
uci set samba.printers.enabled='1'
uci set samba.printers.min_print_space='80000'

uci commit samba
