# etc/uci-defaults/tch_0100-ispconfig
#!/bin/sh

[ -d /etc/ispconfig ] || mkdir /etc/ispconfig

ispconfig_file=/etc/ispconfig/tch_0000-ispconfig

if [ -f $ispconfig_file ]; then
    . $ispconfig_file
fi

exit 0
