# etc/uci-defaults/tch_0020-network-lan
#!/bin/sh

. /lib/functions/uci-defaults.sh

ucidef_set_interface_lan "eth0 eth1 eth2 eth3"

uci set network.lan.ip6assign=64
uci set network.lan.ip6hint=0
uci set network.lan.netmask=255.255.255.0
uci set network.lan.ipaddr=10.1.1.1
uci set network.lan.force_link=1
uci set network.lan.ip6weight=1
uci set network.lan.rpfilter=1

uci set network.VLAN1=interface
uci set network.VLAN1.type=bridge
uci set network.VLAN1.ip6assign=64
uci set network.VLAN1.force_link=1
uci set network.VLAN1.auto=0

uci set network.VLAN2=interface
uci set network.VLAN2.type=bridge
uci set network.VLAN2.ip6assign=64
uci set network.VLAN2.force_link=1
uci set network.VLAN2.auto=0

uci set network.VLAN3=interface
uci set network.VLAN3.type=bridge
uci set network.VLAN3.ip6assign=64
uci set network.VLAN3.force_link=1
uci set network.VLAN3.auto=0

uci set network.Guest1=interface
uci set network.Guest1.proto=static
uci set network.Guest1.ip6assign=64
uci set network.Guest1.netmask=255.255.255.0
uci set network.Guest1.ipaddr=10.1.2.1
uci set network.Guest1.ifname=wl0_1
uci set network.Guest1.force_link=0
uci set network.Guest1.rpfilter=1

uci set network.Guest1_5GHz=interface
uci set network.Guest1_5GHz.proto=static
uci set network.Guest1_5GHz.ip6assign=64
uci set network.Guest1_5GHz.netmask=255.255.255.0
uci set network.Guest1_5GHz.ipaddr=10.1.3.1
uci set network.Guest1_5GHz.ifname=wl1_1
uci set network.Guest1_5GHz.force_link=0
uci set network.Guest1_5GHz.rpfilter=1

ucidef_set_switch "bcmsw" "1" "0"
uci set network.bcmsw.qosimppauseenable=0
uci commit network

exit 0
