# etc/uci-defaults/tch_5001_xmpp
#!/bin/sh

uci set cwmpd.cwmpd_config.xmpp_resource="$(uci get env.var.serial)"
uci commit cwmpd
