# etc/uci-defaults/tch_0020-network-lan
#!/bin/sh

. /lib/functions/uci-defaults.sh

uci set network.lan=interface
uci set network.lan.type=bridge
uci set network.lan.proto=static
uci set network.lan.ipaddr=192.168.1.1
uci set network.lan.netmask=255.255.255.0
uci set network.lan.force_link=0
uci set network.lan.igmp_snooping=1
uci set network.lan.ipv6=0
uci set network.lan.ip6assign=64
uci add_list network.lan.ifname='eth0' 
uci add_list network.lan.ifname='eth1' 
uci add_list network.lan.ifname='eth2' 
uci add_list network.lan.ifname='eth3'
#add config from CONFIG_WLAN_REMOTE_MGMT_ETH_INTERFACE
uci add_list network.lan.ifname='eth5'

ucidef_set_switch "bcmsw" "1" "0"
uci set network.bcmsw.qosimppauseenable=0
uci set network.bcmsw.jumbo=1
uci commit network


uci commit network

exit 0
