# etc/uci-defaults/tch_0021-faultmgmt
#!/bin/sh

uci set dhcp.@dnsmasq[0].logqueries='1'
uci set dhcp.@dnsmasq[0].logfacility='/var/log/dnsmasq'

uci commit dhcp

