# etc/uci-defaults/tch_0071-datausage-web
uci -q del_list web.ruleset_main.rules=datausagemodal
uci add_list web.ruleset_main.rules=datausagemodal

uci -q delete web.datausagemodal
uci set web.datausagemodal=rule
uci set web.datausagemodal.target=/modals/datausage.lp
uci add_list web.datausagemodal.roles=engineer

uci -q del_list web.ruleset_main.rules=datausageajax
uci add_list web.ruleset_main.rules=datausageajax

uci -q delete web.datausageajax
uci set web.datausageajax=rule
uci set web.datausageajax.target=/ajax/datausage.lua
uci add_list web.datausageajax.roles=engineer

uci commit web
