# etc/uci-defaults/tch_0030-network-wan
#!/bin/sh

. /lib/functions/uci-defaults.sh

uci set network.wan=interface
uci set network.wan.proto='pppoe'
uci set network.wan.keepalive=5,10
uci set network.wan.pppoe_backoff=1,32
uci set network.wan.username=`uci get env.var.pppuser`
uci set network.wan.password=`uci get env.var.ppppassword`
uci set network.wan.auto='0'

uci set network.wan6=interface
uci set network.wan6.ifname=@wan
uci set network.wan6.proto=dhcpv6
uci set network.wan6.reqopts='3 12 21 22 23 24 25 31 56 64 67 82 83'
uci set network.wan6.noslaaconly=1
uci set network.wan6.iface_464xlat=0
uci set network.wan6.reqaddress=none

uci set network.ppp_8_35=interface
uci set network.ppp_8_35.ifname='atm_8_35_v0'
uci set network.ppp_8_35.auto='0'

uci set network.ppp_8_81=interface
uci set network.ppp_8_81.ifname='atm_8_81_v0'
uci set network.ppp_8_81.auto='0'

uci set network.atm_8_35_v0=bcmvopi
uci set network.atm_8_35_v0.if='atm_8_35'
uci set network.atm_8_35_v0.vid='0'
uci set network.atm_8_35_v0.defvid='yes'
uci set network.atm_8_35_v0.dscp2pbits='yes'

uci set network.atm_8_81_v0=bcmvopi
uci set network.atm_8_81_v0.if='atm_8_81'
uci set network.atm_8_81_v0.vid='0'
uci set network.atm_8_81_v0.defvid='yes'
uci set network.atm_8_81_v0.dscp2pbits='yes'

uci set network.ptm0_v881=bcmvopi
uci set network.ptm0_v881.if='ptm0'
uci set network.ptm0_v881.vid='881'
uci set network.ptm0_v881.dscp2pbits='yes'

uci commit network

exit 0
