# etc/uci-defaults/tch_0060-transformer
#!/bin/sh
#NG-73564 (removed  uci add_list transformer.@main[0].vendor_patterns='X_TELECOMITALIA_')

if [ ! -f /etc/config/transformer ]; then
	touch /etc/config/transformer
fi

uci show transformer.@main[0] >/dev/null 2>/dev/null
if [ $? -ne 0 ]; then
	uci add transformer main
fi

uci add_list transformer.@main[0].ignore_patterns='^InternetGatewayDevice%.'
uci add_list transformer.@main[0].ignore_patterns='DeviceInfo%.SoftwareBank%.'
uci add_list transformer.@main[0].ignore_patterns='DeviceInfo%.NewSoftwareBank$'
uci add_list transformer.@main[0].unhide_patterns='^Device%.WiFi%.AccessPoint%.{i}%.Security%.WEPKey$'
uci add_list transformer.@main[0].unhide_patterns='^Device%.WiFi%.AccessPoint%.{i}%.Security%.KeyPassphrase$'
uci commit transformer
