# etc/uci-defaults/tch_0060-transformer
#!/bin/sh
#NG-73564 (removed  uci add_list transformer.@main[0].vendor_patterns='X_TELECOMITALIA_')
#NG-85629
#NG-100382 TI: Optical interface data model parameters present in VBNT-K


if [ ! -f /etc/config/transformer ]; then
	touch /etc/config/transformer
fi

uci show transformer.@main[0] >/dev/null 2>/dev/null
if [ $? -ne 0 ]; then
	uci add transformer main
fi

sfp=$(uci get env.rip.sfp)
if [ $sfp -ne "1" ]; then
	uci add_list transformer.@main[0].ignore_patterns='rpc%.optical%.'
	uci add_list transformer.@main[0].ignore_patterns='Device%.Optical%.'
fi

uci add_list transformer.@main[0].vendor_patterns='X_TELECOMITALIA_IT_'

uci commit transformer
