# etc/uci-defaults/tch_0021-firewall-helpers
#!/bin/sh

ln -s /usr/share/fw3/helpers.conf /etc/config/firewall_helpers
