# etc/uci-defaults/tch_5000_versioncusto
#!/bin/sh

********************************************************************
*		Modifiing the SW version to the values required by Customer
*********************************************************************
board=`uci get env.var.hardware_version`


if [ $board == "VBNT-K" ]; then
	uci set versioncusto.override.fwversion_prefix='AGTEF_'
	transl="TI_versions"
elif [ $board == "VBNT-K_SFP" ]; then
	uci set versioncusto.override.fwversion_prefix='AGTHP_'
	transl="TI_versions_SFP"	
elif [ $board == "VBNT-S" ]; then
	uci set versioncusto.override.fwversion_prefix='AGTHP_'
	transl="TI_versions_SFP"	
fi

	
#this will only be used during development/customization and will be removed normally together with closing the build 
#uci set versioncusto.override.fwversion_suffix='ti_open file 16.3.7187-2789001'


#firmware versions conversion table : homeware version -> TI version
#Add new entries in the 'here-document' below (before the EOTX ending line)
#leave the last line : 'Unknown' if no matching version
#top-down precedence order : if actual version does not match first field of first line, try next line
TI_versions(){
cat <<EOTXT
16.3.7172-2789002 1.0.0_001
16.3.7172-2781002 1.0.0_001
16.3.7187-2789001 1.0.0_002
16.3.7187-2781001 1.0.0_002
16.3.7187-2789002 1.0.0_002
16.3.7187-2781002 1.0.0_002
16.3.7187-2789003 1.0.0_002_DUMMY
16.3.7187-2781003 1.0.0_002_DUMMY
16.3.7187-2781007 1.0.0_003_DUMMY
16.3.7196-2789001 1.0.0_003
16.3.7196-2781001 1.0.0_003
16.3.7196-2789002 1.0.0_003
16.3.7196-2781002 1.0.0_003
16.3.7196-2789003 1.0.0_003_DUMMY
16.3.7196-2781003 1.0.0_003_DUMMY
16.3.7271-2789001 1.0.0
16.3.7271-2781001 1.0.0
16.3.7271-2789002 1.0.0
16.3.7271-2781002 1.0.0
16.3.7271-2789003 1.0.0_DUMMY
16.3.7271-2781003 1.0.0_DUMMY
16.3.7318-2789001 1.0.0_005
16.3.7318-2781001 1.0.0_005
16.3.7318-2789002 1.0.0_005
16.3.7318-2781002 1.0.0_005
16.3.7318-2789003 1.0.0_005_DUMMY
16.3.7318-2781003 1.0.0_005_DUMMY
16.3.7343-2789001 1.0.1
16.3.7343-2781001 1.0.1
16.3.7343-2789003 1.0.1
16.3.7343-2781003 1.0.1
16.3.7343-2789004 1.0.1_DUMMY
16.3.7343-2781004 1.0.1_DUMMY
16.3.7391-2789001 1.0.2_001
16.3.7391-2781001 1.0.2_001
16.3.7391-2789002 1.0.2_001
16.3.7391-2781002 1.0.2_001
16.3.7391-2789003 1.0.2_001_DUMMY
16.3.7391-2781003 1.0.2_001_DUMMY
16.3.7498-2789008 1.0.2
16.3.7498-2781008 1.0.2
16.3.7498-2789009 1.0.2
16.3.7498-2781009 1.0.2
16.3.7498-2789010 1.0.2_DUMMY
16.3.7498-2781010 1.0.2_DUMMY
16.3.7573-2789001 1.0.3_001
16.3.7573-2781001 1.0.3_001
16.3.7573-2789002 1.0.3_001
16.3.7573-2781002 1.0.3_001
16.3.7573-2789003 1.0.3_001_DUMMY
16.3.7573-2781003 1.0.3_001_DUMMY
16.3.7636-2789001 1.0.3_002
16.3.7636-2781001 1.0.3_002
16.3.7636-2789002 1.0.3_002
16.3.7636-2781002 1.0.3_002
16.3.7636-2789003 1.0.3_002_DUMMY
16.3.7636-2781003 1.0.3_002_DUMMY
16.3.7636-2789004 1.0.3_003
16.3.7636-2781004 1.0.3_003
16.3.7636-2789005 1.0.3_003
16.3.7636-2781005 1.0.3_003
16.3.7636-2789006 1.0.3_003_DUMMY
16.3.7636-2781006 1.0.3_003_DUMMY
16.3.7636-2789007 1.0.3
16.3.7636-2781007 1.0.3
16.3.7636-2789008 1.0.3
16.3.7636-2781008 1.0.3
16.3.7636-2789009 1.0.3_DUMMY
16.3.7636-2781009 1.0.3_DUMMY
17.1.7711-0009001 1.0.4_002
17.1.7711-0001001 1.0.4_002
17.1.7711-0009002 1.0.4_002
17.1.7711-0001002 1.0.4_002
17.1.7711-0009003 1.0.4_002_DUMMY
17.1.7711-0001003 1.0.4_002_DUMMY
17.1.7745-0009001 1.0.4_003
17.1.7745-0001001 1.0.4_003
17.1.7745-0009010 1.0.4_003
17.1.7745-0001010 1.0.4_003
17.1.7745-0009009 1.0.4_003_DUMMY
17.1.7745-0001009 1.0.4_003_DUMMY
17.1.7765-0009004 1.0.4_004
17.1.7765-0001004 1.0.4_004
17.1.7765-0009014 1.0.4_004
17.1.7765-0001014 1.0.4_004
17.1.7765-0009015 1.0.4_004_DUMMY
17.1.7765-0001015 1.0.4_004_DUMMY
17.1.7765-0009021 1.0.4_005
17.1.7765-0001021 1.0.4_005
17.1.7765-0009022 1.0.4_005
17.1.7765-0001022 1.0.4_005
17.1.7765-0009023 1.0.4_005_DUMMY
17.1.7765-0001023 1.0.4_005_DUMMY
17.1.7841-0009004 1.0.4_006
17.1.7841-0001004 1.0.4_006
17.1.7841-0009005 1.0.4_006
17.1.7841-0001005 1.0.4_006
17.1.7841-0009006 1.0.4_006_DUMMY
17.1.7841-0001006 1.0.4_006_DUMMY
17.1.7882-0009014 1.0.4_007
17.1.7882-0001014 1.0.4_007
17.1.7882-0009015 1.0.4_007
17.1.7882-0001015 1.0.4_007
17.1.7882-0009016 1.0.4_007_DUMMY
17.1.7882-0001016 1.0.4_007_DUMMY
17.1.7882-0009029 1.0.4_009
17.1.7882-0001029 1.0.4_009
17.1.7882-0009030 1.0.4_009
17.1.7882-0001030 1.0.4_009
17.1.7882-0009031 1.0.4_009_DUMMY
17.1.7882-0001031 1.0.4_009_DUMMY
17.3.0165-1689002 1.0.5_001
17.3.0165-1681002 1.0.5_001
17.3.0165-1689005 1.0.5_001
17.3.0165-1681005 1.0.5_001
17.3.0165-1689006 1.0.5_001_DUMMY
17.3.0165-1681006 1.0.5_001_DUMMY
17.3.0177-1689002 1.1.0_001
17.3.0177-1681002 1.1.0_001
17.3.0177-1689004 1.1.0_001
17.3.0177-1681004 1.1.0_001
17.3.0177-1689006 1.1.0_001_DUMMY
17.3.0177-1681006 1.1.0_001_DUMMY
17.3.0177-1689031 1.1.0_002
17.3.0177-1681031 1.1.0_002
17.3.0177-1689032 1.1.0_002
17.3.0177-1681032 1.1.0_002
17.3.0177-1689033 1.1.0_002_DUMMY
17.3.0177-1681033 1.1.0_002_DUMMY
17.3.0177-1689037 1.1.0
17.3.0177-1681037 1.1.0
17.3.0177-1689038 1.1.0
17.3.0177-1681038 1.1.0
17.3.0177-1689039 1.1.0_DUMMY
17.3.0177-1681039 1.1.0_DUMMY
17.3.0238-1689002 1.1.1_001
17.3.0238-1681002 1.1.1_001
17.3.0238-1689006 1.1.1_001
17.3.0238-1681006 1.1.1_001
17.3.0238-1689007 1.1.1_001_DUMMY
17.3.0238-1681007 1.1.1_001_DUMMY
17.3.0268-1689002 1.1.1_002
17.3.0268-1681002 1.1.1_002
17.3.0268-1689006 1.1.1_002
17.3.0268-1681006 1.1.1_002
17.3.0268-1689007 1.1.1_002_DUMMY
17.3.0268-1681007 1.1.1_002_DUMMY
17.3.0268-1689021 1.1.1
17.3.0268-1681021 1.1.1
17.3.0268-1689022 1.1.1
17.3.0268-1681022 1.1.1
17.3.0268-1689023 1.1.1_DUMMY
17.3.0268-1681023 1.1.1_DUMMY
17.3.0307-1689002 1.1.2_001
17.3.0307-1681002 1.1.2_001
17.3.0307-1689004 1.1.2_001
17.3.0307-1681004 1.1.2_001
17.3.0307-1689006 1.1.2_001_DUMMY
17.3.0307-1681006 1.1.2_001_DUMMY
17.3.c.0331-8041004 1.1.2_003
17.3.c.0331-8049004 1.1.2_003
17.3.c.0331-8041005 1.1.2_003
17.3.c.0331-8049005 1.1.2_003
17.3.c.0331-8041006 1.1.2_003_DUMMY
17.3.c.0331-8049006 1.1.2_003_DUMMY
17.3.c.0337-8041014 1.1.2_004
17.3.c.0337-8049014 1.1.2_004
17.3.c.0337-8041015 1.1.2_004
17.3.c.0337-8049015 1.1.2_004
17.3.c.0337-8041016 1.1.2_004_DUMMY
17.3.c.0337-8049016 1.1.2_004_DUMMY
17.3.c.0337-8041018 1.2.0_001
17.3.c.0337-8049018 1.2.0_001
17.3.c.0337-8041023 1.2.0_001
17.3.c.0337-8049023 1.2.0_001
17.3.c.0337-8041024 1.2.0_001_DUMMY
17.3.c.0337-8049024 1.2.0_001_DUMMY
17.3.c.0337-8041029 1.1.2
17.3.c.0337-8049029 1.1.2
17.3.c.0337-8041030 1.1.2
17.3.c.0337-8049030 1.1.2
17.3.c.0337-8041031 1.1.2_DUMMY
17.3.c.0337-8049031 1.1.2_DUMMY
17.3.c.0349-8041002 2.0.0_002
17.3.c.0349-8049002 2.0.0_002
17.3.c.0349-8041006 2.0.0_002
17.3.c.0349-8049006 2.0.0_002
17.3.c.0349-8041007 2.0.0_002_DUMMY
17.3.c.0349-8049007 2.0.0_002_DUMMY
17.3.c.0349-8041009 2.0.0_003
17.3.c.0349-8049009 2.0.0_003
17.3.c.0349-8041011 2.0.0_003
17.3.c.0349-8049011 2.0.0_003
17.3.c.0349-8041013 2.0.0_003_DUMMY
17.3.c.0349-8049013 2.0.0_003_DUMMY
17.3.c.0349-8041018 2.0.0_004
17.3.c.0349-8049018 2.0.0_004
17.3.c.0349-8041019 2.0.0_004
17.3.c.0349-8049019 2.0.0_004
17.3.c.0349-8041020 2.0.0_004_DUMMY
17.3.c.0349-8049020 2.0.0_004_DUMMY
17.3.c.0349-8041024 2.0.0
17.3.c.0349-8049024 2.0.0
17.3.c.0349-8041025 2.0.0
17.3.c.0349-8049025 2.0.0
17.3.c.0349-8041026 2.0.0_DUMMY
17.3.c.0349-8049026 2.0.0_DUMMY
17.3.c.0365-8041002 2.0.1_001
17.3.c.0365-8049002 2.0.1_001
17.3.c.0365-8041005 2.0.1_001
17.3.c.0365-8049005 2.0.1_001
17.3.c.0365-8041006 2.0.1_001_DUMMY
17.3.c.0365-8049006 2.0.1_001_DUMMY
17.3.c.0365-8041010 2.0.1_002
17.3.c.0365-8049010 2.0.1_002
17.3.c.0365-8041012 2.0.1_002
17.3.c.0365-8049012 2.0.1_002
17.3.c.0365-8041013 2.0.1_002_DUMMY
17.3.c.0365-8049013 2.0.1_002_DUMMY
17.3.c.0365-8041017 2.0.1_003
17.3.c.0365-8049017 2.0.1_003
17.3.c.0365-8041018 2.0.1_003
17.3.c.0365-8049018 2.0.1_003
17.3.c.0365-8041019 2.0.1_003_DUMMY
17.3.c.0365-8049019 2.0.1_003_DUMMY
18.3.0313-3161004 2.1.0_001
18.3.0313-3169004 2.1.0_001
18.3.0313-3161005 2.1.0_001
18.3.0313-3169005 2.1.0_001
18.3.0313-3161006 2.1.0_001_DUMMY
18.3.0313-3169006 2.1.0_001_DUMMY
18.3.0335-3161004 2.1.0_002
18.3.0335-3169004 2.1.0_002
18.3.0335-3161005 2.1.0_002
18.3.0335-3169005 2.1.0_002
18.3.0335-3161006 2.1.0_002_DUMMY
18.3.0335-3169006 2.1.0_002_DUMMY
.*	1.99.99.99 
EOTXT
}



#perform a translation (table based in this example : output 2nd field of matching version entry in PXS_versions table above)
#input params:
#   none
TI_translate_swversions(){
	tver=$($transl | awk -v srch="$(uci get env.var.friendly_sw_version_activebank)" 'srch ~ $1 {print $2}' | head -1)
	uci set env.var.friendly_sw_version_activebank=`uci get versioncusto.override.fwversion_prefix`"$tver"
	uci set versioncusto.override.fwversion_override="$tver"
#Populate user-agent
#user-agent format for TI: vendorname / <HWversion>/ <FWversion> / <SWversion>
	uci set mmpbxrvsipnet.sip_net.user_agent="$(uci get env.var.company_name) / $(uci get env.var.hardware_version) / `uci get versioncusto.override.fwversion_prefix`"$tver" / `uci get versioncusto.override.fwversion_prefix`"$tver""
	uci commit mmpbxrvsipnet
	
	tver=$($transl | awk -v srch="$(uci get env.var.friendly_sw_version_passivebank)" 'srch ~ $1 {print $2}' | head -1)
	uci set env.var.friendly_sw_version_passivebank=`uci get versioncusto.override.fwversion_prefix`"$tver"
#if no conversion needed : comment previous 3 lines
uci commit env
}



TI_translate_swversions



SW_Version=`lua -e "t=require'transformer.shared.banktable' value=t.getCurrentVersion() print(string.match(value,'([^%-]+)%-'))"`

#this will only be used during development/customization and will be removed normally together with closing the build 
#uci set versioncusto.override.fwversion_suffix='_Custo005'


uci commit versioncusto
