# etc/uci-defaults/tch_0071-mobiled-web
uci -q del_list web.ruleset_main.rules=ltemodal
uci add_list web.ruleset_main.rules=ltemodal

uci -q delete web.ltemodal
uci set web.ltemodal=rule
uci set web.ltemodal.target=/modals/lte-modal.lp
uci add_list web.ltemodal.roles=admin
uci add_list web.ltemodal.roles=engineer

uci -q del_list web.ruleset_main.rules=ltedoctormodal
uci add_list web.ruleset_main.rules=ltedoctormodal

uci -q delete web.ltedoctormodal
uci set web.ltedoctormodal=rule
uci set web.ltedoctormodal.target=/modals/lte-doctor.lp
uci add_list web.ltedoctormodal.roles=admin
uci add_list web.ltedoctormodal.roles=engineer

uci -q del_list web.ruleset_main.rules=ltesms
uci add_list web.ruleset_main.rules=ltesms

uci -q delete web.ltesms
uci set web.ltesms=rule
uci set web.ltesms.target=/modals/lte-sms.lp
uci add_list web.ltesms.roles=admin
uci add_list web.ltesms.roles=engineer

uci -q del_list web.ruleset_main.rules=lteprofiles
uci add_list web.ruleset_main.rules=lteprofiles

uci -q delete web.lteprofiles
uci set web.lteprofiles=rule
uci set web.lteprofiles.target=/modals/lte-profiles.lp
uci add_list web.lteprofiles.roles=admin
uci add_list web.lteprofiles.roles=engineer

uci -q del_list web.ruleset_main.rules=lteradioparameters
uci add_list web.ruleset_main.rules=lteradioparameters

uci -q delete web.lteradioparameters
uci set web.lteradioparameters=rule
uci set web.lteradioparameters.target=/ajax/radioparameters.lua
uci add_list web.lteradioparameters.roles=admin
uci add_list web.lteradioparameters.roles=engineer

uci -q del_list web.ruleset_main.rules=lteajaxsms
uci add_list web.ruleset_main.rules=lteajaxsms

uci -q delete web.lteajaxsms
uci set web.lteajaxsms=rule
uci set web.lteajaxsms.target=/ajax/sms.lua
#NG-93140 uci add_list web.lteajaxsms.roles=admin
uci add_list web.lteajaxsms.roles=engineer

uci -q del_list web.ruleset_main.rules=lteajaxmobiletab
uci add_list web.ruleset_main.rules=lteajaxmobiletab

uci -q delete web.lteajaxmobiletab
uci set web.lteajaxmobiletab=rule
uci set web.lteajaxmobiletab.target=/ajax/mobiletab.lua
uci add_list web.lteajaxmobiletab.roles=guest
uci add_list web.lteajaxmobiletab.roles=admin
uci add_list web.lteajaxmobiletab.roles=engineer

uci -q del_list web.ruleset_main.rules=ltenetworkscan
uci add_list web.ruleset_main.rules=ltenetworkscan

uci -q delete web.ltenetworkscan
uci set web.ltenetworkscan=rule
uci set web.ltenetworkscan.target=/ajax/networkscan.lua
uci add_list web.ltenetworkscan.roles=admin
uci add_list web.ltenetworkscan.roles=engineer

uci -q del_list web.ruleset_main.rules=lteajaxsim
uci add_list web.ruleset_main.rules=lteajaxsim

uci -q delete web.lteajaxsim
uci set web.lteajaxsim=rule
uci set web.lteajaxsim.target=/ajax/sim.lua
uci add_list web.lteajaxsim.roles=admin
uci add_list web.lteajaxsim.roles=engineer

uci -q del_list web.ruleset_main.rules=ltesim
uci add_list web.ruleset_main.rules=ltesim

uci -q delete web.ltesim
uci set web.ltesim=rule
uci set web.ltesim.target=/modals/lte-sim.lp
uci add_list web.ltesim.roles=admin
uci add_list web.ltesim.roles=engineer

uci -q del_list web.ruleset_main.rules=ltedoctorajax
uci add_list web.ruleset_main.rules=ltedoctorajax

uci -q delete web.ltedoctorajax
uci set web.ltedoctorajax=rule
uci set web.ltedoctorajax.target=/ajax/lte-doctor.lua
uci add_list web.ltedoctorajax.roles=admin
uci add_list web.ltedoctorajax.roles=engineer

uci -q del_list web.ruleset_main.rules=ltedoctor
uci add_list web.ruleset_main.rules=ltedoctor

uci -q delete web.ltedoctor
uci set web.ltedoctor=rule
uci set web.ltedoctor.target=/lte-doctor.lua
uci add_list web.ltedoctor.roles=admin
uci add_list web.ltedoctor.roles=engineer

uci commit web
