# etc/uci-defaults/tch_0036-pppoe-relay
#!/bin/sh

uci set network.lan.pppoerelay="$(uci get network.wan.ifname)"

uci commit network

exit 0
