# etc/uci-defaults/tch_0090-remove
file=/etc/ssl/certs/653b494a.0
if [ -f $file ] ; then
    mkdir -p /etc/ssl/acs-cert/
    cp $file /etc/ssl/acs-cert/
fi

file=/etc/ssl/certs/AGprivkey.key
if [ -f $file ] ; then
    mkdir -p /etc/ssl/acs-cert/
    cp $file /etc/ssl/acs-cert/
fi

file=/etc/ssl/certs/BaltimoreCybertrustRoot.cer
if [ -f $file ] ; then
    mkdir -p /etc/ssl/acs-cert/
    cp $file /etc/ssl/acs-cert/
fi

file=/etc/ssl/certs/d6d84089.0
if [ -f $file ] ; then
    mkdir -p /etc/ssl/acs-cert/
    cp $file /etc/ssl/acs-cert/
fi

file=/etc/ssl/certs/rootTLC-CA.cer
if [ -f $file ] ; then
    mkdir -p /etc/ssl/acs-cert/
    cp $file /etc/ssl/acs-cert/
fi
file=/etc/ssl/certs/4042bcee.0
if [ -f $file ] ; then
    rm $file
fi

file=/etc/ssl/certs/6d41d539.0
if [ -f $file ] ; then
    rm $file
fi

file=/etc/ssl/certs/8cb5ee0f.0
if [ -f $file ] ; then
    rm $file
fi

file=/etc/ssl/certs/aee5f10d.0
if [ -f $file ] ; then
    rm $file
fi

file=/etc/ssl/certs/de6d66f3.0
if [ -f $file ] ; then
    rm $file
fi

file=/etc/ssl/certs/ce5e74ef.0
if [ -f $file ] ; then
    rm $file
fi

file=/www/cards/016_DMZ.lp
if [ -f $file ] ; then
    rm $file
fi

file=/www/cards/017_dmz_mac.lp
if [ -f $file ] ; then
    rm $file
fi

file=/www/cards/018_eth_speed.lp
if [ -f $file ] ; then
    rm $file
fi

file=/www/cards/099_smartHome.lp
if [ -f $file ] ; then
    rm $file
fi

file=/www/cards/snippets/003_internet_default.lp
if [ -f $file ] ; then
    rm $file
fi

file=/www/cards/snippets/003_ipoe_pvc_vlan.lp
if [ -f $file ] ; then
    rm $file
fi

file=/www/cards/snippets/003_internet_wan_sensing.lp
if [ -f $file ] ; then
    rm $file
fi

file=/www/cards/snippets/003_internet_pvc_vlan.lp
if [ -f $file ] ; then
    rm $file
fi

file=/www/cards/snippets/003_internet_ipoe_routed.lp
if [ -f $file ] ; then
    rm $file
fi

file=/www/cards/095_snmp.lp
if [ -f $file ] ; then
    rm $file
fi

file=/www/docroot/modals/diagnostics-igmpproxy-modal.lp
if [ -f $file ] ; then
    rm $file
fi

file=/www/docroot/modals/ipoe-modal.lp
if [ -f $file ] ; then
    rm $file
fi

file=/www/docroot/modals/dmz-modal.lp
if [ -f $file ] ; then
    rm $file
fi

file=/www/docroot/modals/dmz-macmodal.lp
if [ -f $file ] ; then
    rm $file
fi

file=/www/docroot/modals/ethspeed-modal.lp
if [ -f $file ] ; then
    rm $file
fi

file=/www/docroot/modals/bridge-modal.lp
if [ -f $file ] ; then
    rm $file
fi

file=/www/docroot/modals/snmpd-agent-modal.lp
if [ -f $file ] ; then
    rm $file
fi

file=/www/docroot/modals/snmpd-com2sec-modal.lp
if [ -f $file ] ; then
    rm $file
fi

file=/www/docroot/modals/snmpd-global-modal.lp
if [ -f $file ] ; then
    rm $file
fi

file=/www/docroot/modals/snmpd-system-modal.lp
if [ -f $file ] ; then
    rm $file
fi

file=/www/docroot/modals/smarthome-modal.lp
if [ -f $file ] ; then
    rm $file
fi
