# etc/uci-defaults/tch_0060-transformer
#!/bin/sh

if [ ! -f /etc/config/transformer ]; then
	touch /etc/config/transformer
fi

uci show transformer.@main[0] >/dev/null 2>/dev/null
if [ $? -ne 0 ]; then
	uci add transformer main
fi

uci add_list transformer.@main[0].ignore_patterns='^InternetGatewayDevice%.'
uci add_list transformer.@main[0].ignore_patterns='DeviceInfo%.SoftwareBank%.'
uci add_list transformer.@main[0].ignore_patterns='DeviceInfo%.NewSoftwareBank$'
uci add_list transformer.@main[0].ignore_patterns='X_000E50_EnableIPv6$'
uci add_list transformer.@main[0].ignore_patterns='Device.DHCPv4.Server%.'
uci add_list transformer.@main[0].ignore_patterns='Device.DHCPv6.Server%.'
uci add_list transformer.@main[0].ignore_patterns='Device.DHCPv6.Client%.{i}%.Renew$'
uci add_list transformer.@main[0].ignore_patterns='Device.DHCPv6.Client%.{i}%.RequestPrefixes$'
uci add_list transformer.@main[0].ignore_patterns='Device.DHCPv6.Client%.{i}%.DUID$'
uci add_list transformer.@main[0].ignore_patterns='Device.DHCPv6.Client%.{i}%.Enable$'
uci add_list transformer.@main[0].ignore_patterns='Device.DHCPv6.Client%.{i}%.X_000E50_ReleaseRenew$'
uci add_list transformer.@main[0].ignore_patterns='Device.DHCPv6.Client%.{i}%.RequestedOptions$'
uci add_list transformer.@main[0].ignore_patterns='Device.DHCPv6.Client%.{i}%.RequestAddresses$'
uci add_list transformer.@main[0].ignore_patterns='Device.DHCPv6.Client%.{i}%.SupportedOptions$'
uci add_list transformer.@main[0].ignore_patterns='Device.DHCPv6.Client%.{i}%.Status$'
uci add_list transformer.@main[0].ignore_patterns='Device.DHCPv6.Client%.{i}%.RapidCommit$'
uci add_list transformer.@main[0].ignore_patterns='Device.DHCPv6.Client%.{i}%.Alias$'
uci add_list transformer.@main[0].ignore_patterns='Device.DHCPv6.Client%.{i}%.Interface$'
uci add_list transformer.@main[0].ignore_patterns='Device.DHCPv6.Client%.{i}%.SentOption%.{i}%.'
uci add_list transformer.@main[0].ignore_patterns='Device%.DeviceInfo%.NetworkProperties%.'
uci add_list transformer.@main[0].ignore_patterns='Device.DeviceInfo%.X_000E50_LED%.PowerLEDState$'
uci add_list transformer.@main[0].ignore_patterns='Device.DeviceInfo%.X_000E50_LED%.WirelessLEDState_5G$'
uci add_list transformer.@main[0].ignore_patterns='Device.DeviceInfo%.X_000E50_LED%.WirelessLEDState$'
uci add_list transformer.@main[0].ignore_patterns='Device.DeviceInfo%.X_000E50_LED%.IPTVLED$'
uci add_list transformer.@main[0].ignore_patterns='Device.DeviceInfo%.X_000E50_LED%.DSLLEDState$'
uci add_list transformer.@main[0].ignore_patterns='Device.DeviceInfo%.X_000E50_LED%.EthernetLED$'
uci add_list transformer.@main[0].ignore_patterns='Device.DeviceInfo%.X_000E50_LED%.WPSLED$'
uci add_list transformer.@main[0].ignore_patterns='Device%.DynamicDNS%.'
uci add_list transformer.@main[0].ignore_patterns='Device%.Firewall%.'
uci add_list transformer.@main[0].ignore_patterns='Device%.GRE%.'
uci add_list transformer.@main[0].ignore_patterns='Device%.HPNA%.'
uci add_list transformer.@main[0].ignore_patterns='Device%.LANConfigSecurity%.'
uci add_list transformer.@main[0].ignore_patterns='Device%.Leds%.'
uci add_list transformer.@main[0].ignore_patterns='Device%.ManagementServer%.ManageableDevice%.'
uci add_list transformer.@main[0].ignore_patterns='Device%.MQTT%.'
uci add_list transformer.@main[0].ignore_patterns='Device%.NAT%.'
uci add_list transformer.@main[0].ignore_patterns='Device%.NeighborDiscovery%.'
uci add_list transformer.@main[0].ignore_patterns='Device%.PPP%.'
uci add_list transformer.@main[0].ignore_patterns='Device%.RouterAdvertisement%.'
uci add_list transformer.@main[0].ignore_patterns='Device%.Routing%.'
uci add_list transformer.@main[0].ignore_patterns='Device%.SoftwareModules%.'
uci add_list transformer.@main[0].ignore_patterns='Device%.StorageService%.'
uci add_list transformer.@main[0].ignore_patterns='Device%.UserInterface%.'
uci add_list transformer.@main[0].ignore_patterns='Device%.UPnP%.'
uci add_list transformer.@main[0].ignore_patterns='Device%.WiFi%.MultiAP%.'
uci add_list transformer.@main[0].ignore_patterns='Device%.X_000E50_AccessControl%.'
uci add_list transformer.@main[0].ignore_patterns='Device%.X_000E50_DNSRule%.'
uci add_list transformer.@main[0].ignore_patterns='Device%.X_000E50_GUILastAccess%.'
uci add_list transformer.@main[0].ignore_patterns='Device%.X_000E50_IPSet%.'
uci add_list transformer.@main[0].ignore_patterns='Device%.X_000E50_Iperf%.'
uci add_list transformer.@main[0].ignore_patterns='Device%.X_000E50_L2TPIPSecServer%.'
uci add_list transformer.@main[0].ignore_patterns='Device%.X_000E50_MWANHost%.'
uci add_list transformer.@main[0].ignore_patterns='Device%.X_000E50_MWANPolicy%.'
uci add_list transformer.@main[0].ignore_patterns='Device%.X_000E50_MWANRule%.'
uci add_list transformer.@main[0].ignore_patterns='Device%.X_000E50_MultiAP%.'
uci add_list transformer.@main[0].ignore_patterns='Device%.X_000E50_TimeOfDay%.'
uci add_list transformer.@main[0].ignore_patterns='Device%.Services%.X_000E50_ALG%.'
uci add_list transformer.@main[0].ignore_patterns='Device%.Services%.X_000E50_DNSProxy%.'
uci add_list transformer.@main[0].ignore_patterns='Device%.Services%.X_000E50_HotspotService%.'
uci add_list transformer.@main[0].ignore_patterns='Device%.Services%.X_000E50_IGMPProxy%.'
uci add_list transformer.@main[0].ignore_patterns='Device%.Services%.X_000E50_Internet%.'
uci add_list transformer.@main[0].ignore_patterns='Device%.Services%.X_000E50_RemoteAccess%.'
uci add_list transformer.@main[0].ignore_patterns='Device.IP.ActivePort%.{i}%.LocalIPAddress$'
uci add_list transformer.@main[0].ignore_patterns='Device.IP%.Diagnostics%.X_000E50_BFDEcho%.'
uci add_list transformer.@main[0].ignore_patterns='Device.Services%.X_000E50_AccessManagement%.SIP%.AllowIP$'

uci commit transformer
