# etc/uci-defaults/tch_0005-conloglevel
#!/bin/sh

# Disable console logging
uci set system.@system[0].conloglevel=1
