# etc/uci-defaults/tch_5000-custovar
uci set env.custovar.ONTmode='0'
uci set env.custovar.ONTmodeGUIbutton='0'

uci commit env
