# etc/uci-defaults/tch_0020-network-lan
#!/bin/sh

. /lib/functions/uci-defaults.sh

ucidef_set_interface_lan "eth0 eth1 eth2 eth3"

uci set network.lan.ipaddr=192.168.1.254
uci set network.lan.ip6assign=64
uci set network.lan.force_link=0
uci set network.lan.ipv6=1

ucidef_set_switch "bcmsw" "1" "0"
uci set network.bcmsw.qosimppauseenable=0
uci set network.bcmsw.jumbo=0
uci commit network


#add config from CONFIG_WLAN_REMOTE_MGMT_ETH_INTERFACE

uci set network.lan.ifname="`uci -q get \"network.lan.ifname\"` eth5"

uci commit network

exit 0
