# etc/uci-defaults/tch_0030-network-wan
#!/bin/sh

. /lib/functions/uci-defaults.sh

uci set network.wan=interface
uci set network.wan.proto='dhcp'
uci set network.wan.reqopts='1 3 6 15 33 42 43 51 121 249 212'
uci set network.wan.vendorid=`uci get env.var.dhcp_option60`
uci set network.wan.authfail='0'
uci set network.wan.iface6rd='6rd'
uci set network.wan.zone6rd='wan'
uci set network.wan.mtu='1500'
uci set network.wan.initboot='1'
uci set network.wan.ipv6='0'
uci set network.wan.auto='0'
uci set network.wan.notused=`uci get env.rip.wifi_mac`

uci set network.mgmt=interface
uci set network.mgmt.vendorid=`uci get env.var.dhcp_option60`

uci set network.vlan_ptm0=device
uci set network.vlan_ptm0.type='8021q'
uci set network.vlan_ptm0.ifname='ptm0'
uci set network.vlan_ptm0.mtu='1500'
uci set network.vlan_ptm0.name='vlan_ptm0'
uci set network.vlan_ptm0.vid='100'

uci set network.wantag=interface
uci set network.wantag.proto='dhcp'
uci set network.wantag.reqopts='1 3 6 15 33 42 51 121 249 212'
uci set network.wantag.vendorid=`uci get env.var.dhcp_option60`
uci set network.wantag.authfail='0'
uci set network.wantag.iface6rd='6rd'
uci set network.wantag.zone6rd='wan'
uci set network.wantag.mtu='1500'
uci set network.wantag.initboot='1'
uci set network.wantag.ipv6='0'
uci set network.wantag.auto='0'
uci set env.custovar.wan_tag='0'
uci set network.wwan=interface
uci set network.wwan.auto='1'

uci set network.vlan_wan=interface
uci set network.vlan_wan.ifname='vlan_ptm0'

uci set network.if0atm=device
uci set network.if0atm.enable='1'
uci set network.if0atm.name='atm0'
uci set network.if0atm.macaddr=`uci get env.rip.eth_mac`

uci set network.if1atm=device
uci set network.if1atm.enable='1'
uci set network.if1atm.name='atm1'
uci set network.if1atm.macaddr=`uci get env.rip.wifi_mac`

uci set network.vlan_eth4=device
uci set network.vlan_eth4.type='8021q'
uci set network.vlan_eth4.ifname='eth4'
uci set network.vlan_eth4.mtu='1500'
uci set network.vlan_eth4.name='vlan_eth4'
uci set network.vlan_eth4.vid='835'

uci set network.vlan_wan_eth4=interface
uci set network.vlan_wan_eth4.ifname='vlan_eth4'

#############################################
# Sfp module support two default management IP:
# untag "192.168.2.1" and tag "192.168.10.1",
# each of them could be used for telnet and
# management.
# For connect sfp module, board need match at least
# one of the two management IP.
# Now generic will add two interfaces "sfp" and "sfptag"
# to match the two management IP.
# "sfp"--> SFP untag IP "192.168.2.1"
# "sfptag"--> SFP tag IP "192.168.10.1"
############################################
sfp=`uci get env.rip.sfp`

if [ "$sfp" = "1" ] ; then
  uci set network.sfp=interface
  uci set network.sfp.proto='static'
  uci set network.sfp.ifname='eth4'
  uci set network.sfp.ipaddr='169.0.0.2'
  uci set network.sfp.netmask='255.255.255.0'
fi

uci commit network

uci show dmordering >/dev/null 2>/dev/null || touch /etc/config/dmordering

uci set dmordering.IPInterface=ordering
uci add_list dmordering.IPInterface.order=loopback
uci add_list dmordering.IPInterface.order=lan
uci add_list dmordering.IPInterface.order=Guest
uci add_list dmordering.IPInterface.order=sfp
uci add_list dmordering.IPInterface.order=wan
uci add_list dmordering.IPInterface.order=wwan
uci add_list dmordering.IPInterface.order=wantag
uci add_list dmordering.IPInterface.order=mgmt

uci set dmordering.dhcpv4=ordering
uci set dmordering.dhcpv4.name=Device.DHCPv4.Server.Pool.{i}.
uci add_list dmordering.dhcpv4.order=lan
uci add_list dmordering.dhcpv4.order=Guest

uci commit env
uci commit dmordering

exit 0
