# etc/uci-defaults/tch_0020-network-lan
#!/bin/sh

. /lib/functions/uci-defaults.sh

#all eth ports in the bridge
ucidef_set_interface_lan "eth0 eth1 eth2 eth3"

#set br-lan to 192.168.1.254
uci set network.lan.ipaddr=192.168.1.254
uci set network.lan.ip6assign=64
uci set network.lan.ip6ifaceid='eui64'
uci set network.lan.force_link=0
uci set network.lan.ipv6=1

ucidef_set_switch "bcmsw" "1" "0"
uci set network.bcmsw.qosimppauseenable=1
uci commit network

exit 0
