# etc/uci-defaults/tch_0072-web-fastweb
uci -q delete web.speedlog
uci set web.speedlog=connectionmgr
uci set web.speedlog.auto=0
uci set web.speedlog.time=0000

touch /etc/config/fastweb
uci set fastweb.wan=api
uci set fastweb.wan.port=443
uci set fastweb.wan.auth_required=1

uci set fastweb.lan=api
uci set fastweb.lan.port=80

uci set fastweb.port_shaping=portshaping
uci set fastweb.port_shaping.total_download_bw=0
uci set fastweb.port_shaping.ds_disable=1

uci set fastweb.webui=webui
uci add_list fastweb.webui.urls='1|1|'
uci add_list fastweb.webui.urls='2|1|'
uci add_list fastweb.webui.urls='3|1|'
uci add_list fastweb.webui.urls='4|1|'
uci add_list fastweb.webui.urls='5|1|'
uci add_list fastweb.webui.urls='6|1|'
uci add_list fastweb.webui.urls='7|1|'
uci add_list fastweb.webui.urls='8|1|'
uci add_list fastweb.webui.urls='9|1|'
uci add_list fastweb.webui.urls='10|1|'
uci add_list fastweb.webui.urls='11|1|'
uci add_list fastweb.webui.urls='12|1|'
uci add_list fastweb.webui.urls='13|1|'
uci add_list fastweb.webui.urls='14|1|'
uci add_list fastweb.webui.urls='15|1|'
uci add_list fastweb.webui.urls='16|1|'
uci add_list fastweb.webui.urls='17|1|'
uci add_list fastweb.webui.urls='18|1|'
uci add_list fastweb.webui.urls='19|1|'
uci add_list fastweb.webui.urls='20|1|'

uci commit
