# etc/uci-defaults/tch_0072-web-fastweb
uci -q delete web.speedlog
uci set web.speedlog=connectionmgr
uci set web.speedlog.auto=0
uci set web.speedlog.time=9999

touch /etc/config/fastweb
uci set fastweb.wan=api
uci set fastweb.wan.port=443
uci set fastweb.wan.auth_required=1

uci set fastweb.lan=api
uci set fastweb.lan.port=80

uci set fastweb.port_shaping=portshaping
uci set fastweb.port_shaping.total_download_bw=0
uci set fastweb.port_shaping.ds_disable=1

uci set fastweb.webui=fw_settings
uci set fastweb.webui.key_00=0
uci set fastweb.webui.key_01=0
uci set fastweb.webui.parental_control_mode=0

uci commit
