# etc/uci-defaults/tch_0005-conloglevel
#!/bin/sh

uci set system.@system[0].conloglevel=1
uci commit system
