# etc/uci-defaults/tch_9001-upgrade_sts
#!/bin/sh
uci -q get mobiled.globals.pdn_retry_timer && exit 0
uci set mobiled.globals.pdn_retry_timer=720
uci delete mobiled.@mobiled_state[12].entryexits
uci set mobiled.@mobiled_state[13].entryexits=DataSessionSetup-EntryExit
uci set mobiled.@mobiled_state[11].timeout=60
uci commit mobiled
uci set tod.global.tod_enabled=1
uci set tod.global.tracelevel=4
uci set tod.global.time_change_monfreq=10
uci commit tod
exit 0
