# etc/uci-defaults/tch_9001-upgrade_sts
#!/bin/sh
uci set network.lan.force_link=1
uci commit network
uci delete hotspotd.TLS5G
uci commit hotspotd
exit 0
