# etc/uci-defaults/tch_0030-network-wan
#!/bin/sh

. /lib/functions/uci-defaults.sh

uci set network.wan=interface
uci set network.wan.proto=dhcp
uci set network.wan.ifname=ptm0
uci set network.wan.reqopts='1 3 6 15 33 42 43 51 121 249'
uci set network.wan.vendorid=`uci get env.var.dhcp_option60`
uci set network.wan.authfail='0'
uci set network.wan.iface6rd='6rd'
uci set network.wan.zone6rd='wan'
uci set network.wan.mtu='1500'
uci set network.wan.initboot='1'
uci set network.wan.ipv6='0'

uci set network.wwan=interface
uci set network.wwan.auto='0'


uci commit network

exit 0
