# etc/uci-defaults/tch_9001-upgrade_sts
#!/bin/sh
uci set wireless.ap0.pmf='enabled'
uci set wireless.ap1.pmf='enabled'
uci set wireless.ap2.pmf='enabled'
uci set wireless.ap3.pmf='enabled'
uci commit wireless
exit 0
