# etc/uci-defaults/tch_0031-pppoe-relay
#!/bin/sh

uci add_list network.lan.pppoerelay="wanptm0"
uci commit network

exit 0
