# srv/lxc/lxc_ee/rootfs/etc/init.d/installer
#!/bin/sh /etc/rc.common

START=05

boot() {
	for ipk in $(ls -1 /etc/ipks/*.ipk); do
		opkg install ${ipk} && rm ${ipk}
	done
}

