# srv/lxc/lxc_ee/rootfs/etc/init.d/lxc_start
#!/bin/sh /etc/rc.common

START=05

boot() {
  rm -f /etc/.booted
}
