# srv/lxc/lxc_ee/rootfs/etc/init.d/zend
#!/bin/sh /etc/rc.common

START=99
STOP=01

# This is tze end, my only friend, z end.
boot() {
  touch /etc/.booted
}

stop() {
  rm -f /etc/.booted
}
