# etc/hotplug.d/dhcp/01-dumaos-dhcp-event
#!/bin/sh

/dumaos/apps/system/com.netdumasoftware.neighwatch/dhcp-event.lua "$ACTION" "$MACADDR" "$IPADDR" "$HOSTNAME"
