# etc/init.d/brcm_kthrd_prios
#!/bin/sh /etc/rc.common

. $IPKG_INSTROOT/lib/functions/functions-tch.sh

# After modules loaded and XDSL started
START=19

# Match the priorities Broadcom set on several kthreads.
# Information based on eid_bcm_kthreads.txt, wlmngr and running threads.

start() {
    # USB workqueue thread
    setrtprio khubd 5

    # Broadcom Buffer Pool Manager thread
    setrtprio bpm 75

    # Broadcom FAP GSO LOOPBACK thread
    setrtprio fapGsoLoopBk 5

    # Broadcom bcmenet driver threads
    setrtprio bcmsw_rx 5
    setcpumask bcmsw_rx 1

    ## WLAN ##
    # Thread prio
    # WFD < 4.16L03
    setrtprio wfd 5

    # WFD >= 4.16L03
    setrtprio wfd0-thrd 5
    setrtprio wfd1-thrd 5
    setrtprio wfd2-thrd 5

    # Wlan driver
    setrtprio wl0-kthrd 5
    setrtprio wl1-kthrd 5
    setrtprio dhd0_dpc 5
    setrtprio dhd1_dpc 5
    setrtprio dhd2_dpc 5

    # CPU affinity (see logic in userspace/private/apps/wlan/wlmngr.c)
    cpu_type=`uname -m`

    nr_wlan=1
    if [ "$(process_exists wl1-kthrd)" == "1" ] || [ "$(process_exists dhd1_dpc)" == "1" ] ; then
        nr_wlan=2
    fi
    # Special board (GBNT-Q family) has 3 Full dongle wifi chipset
    if [ "$(process_exists dhd2_dpc)" == "1" ] ; then
        nr_wlan=3
    fi

    if [ "$nr_wlan" == "1" ] ; then
       #If arm: bind on TP0 + also bind wfd
       #If mips: bind wl on TP1, dhd on TP0

       if [ "$cpu_type" == "armv7l" ] ; then
           setcpumask wl0-kthrd 1
           setcpumask dhd0_dpc  1
           setcpumask wfd0-thrd 1
       else
           setcpumask wl0-kthrd 2
           setcpumask dhd0_dpc  1
       fi
    elif [ "$nr_wlan" == "2" ] ; then
      #If wl0 and dhd1: put dhd1 on tp0 to line up with brcm ref where detection order is different
      #If wl0 and wl1: also fix skbFreeTask to TP0 (to improve 11ac perf)
      #If arm: also bind wfd
      if [ "$(process_exists wl0-kthrd)" == "1" ] && [ "$(process_exists dhd1_dpc)" == "1" ] ; then
        setcpumask wl0-kthrd 2
        setcpumask dhd1_dpc 1
        if [ "$cpu_type" == "armv7l" ] || [ "$cpu_type" == "aarch64" ] ; then
           setcpumask wfd0-thrd 2
           setcpumask wfd1-thrd 1
        fi
      else
        setcpumask wl0-kthrd 1
        setcpumask dhd0_dpc  1
        setcpumask wl1-kthrd 2
        setcpumask dhd1_dpc  2
        setcpumask skbFreeTask   1 #Before 5.02L03
        setcpumask skb_free_task 1 #From   5.02L03
        if [ "$cpu_type" == "armv7l" ] ; then
           setcpumask wfd0-thrd 1
           setcpumask wfd1-thrd 2
        fi
      fi
    else
        # Special board (GBNT-Q family) has 3 Full dongle wifi chipset
        setcpumask dhd0_dpc 1
        setcpumask dhd1_dpc 2
        setcpumask dhd2_dpc 4
        setcpumask wfd0-thrd 1
        setcpumask wfd1-thrd 2
        setcpumask wfd2-thrd 4
    fi

    #Broadcom traffic manager kernel thread affinity(see eid_bcm_kthreads.txt)
    if [ "$(process_exists bcm_tm_thread)" == "1" ] ; then
	setcpumask bcm_tm_thread 2
	setrtprio bcm_tm_thread 99 fifo
    fi
    #############################################
    #Configure priority for Technicolor Packages#
    #############################################

}
