# etc/init.d/gpon-uni
#!/bin/sh /etc/rc.common

START=19

. $IPKG_INSTROOT/lib/network/gpon.sh


start() {
	gpon_load_uni
}

stop() {
	gpon_unload_uni
}
