# etc/init.d/mwan
#!/bin/sh /etc/rc.common

START=21

start()
{
	$IPKG_INSTROOT/usr/lib/mwan/setup.sh start
}

stop()
{
	$IPKG_INSTROOT/usr/lib/mwan/setup.sh stop
}
