# etc/init.d/random_seed
#!/bin/sh /etc/rc.common
# Copyright (C) 2015 Technicolor Delivery Technologies, SAS

STOP=99

stop() {
	random_seed "stop"
}

