# etc/init.d/sfpmon
#!/bin/sh /etc/rc.common


START=99

#service for sfp monitor currently only do insert sff8472_i2c ko, will do more in future

start() {
	insmod sff8472_i2c
}


