# etc/init.d/wifi-doctor-agent
#!/bin/sh /etc/rc.common

START=99
SERVICE_DAEMONIZE=1

start() {
        config_load wifi_doctor_agent
        config_get_bool enabled config enabled 0
        [ ${enabled} -gt 0 ] || return 1
        logger -t WIFIDR "Starting wifi-doctor-agent"
        service_start /usr/bin/wifi-doctor-agent
}

stop() {
        logger -t WIFIDR "Stopping wifi-doctor-agent"
        service_stop /usr/bin/wifi-doctor-agent
}

reload() {
        logger -t WIFIDR "Restarting wifi-doctor-agent"
        restart
}
