# etc/hotplug.d/net/01-qos
#!/bin/sh

[ "${ACTION}" = add -a "${DEVTYPE}" = "vlan" ] || exit 0

. /usr/lib/qos/functions.sh

# basic check if QoS started
qos_started || exit 0

/usr/bin/qos -q vlan "${DEVICENAME}"
