# etc/hotplug.d/net/10_ethoam
#!/bin/sh

if [ "$ACTION" = add ]; then
	/etc/init.d/ethoam ifadd_event $DEVICENAME
fi
