# usr/lib/opkg/info/airiq-controller.list
/etc/init.d/airiq-controller
/usr/bin/airiq-controller
