# usr/lib/opkg/info/base-files.list
/lib/functions/syslog.sh
/etc/openwrt_version
/etc/hotplug.d/usb/12-tch_usb
/etc/resolv.conf
/etc/init.d/done
/bin/config_generate
/var
/etc/init.d/bcm_modules
/lib/preinit/00_preinit.conf
/etc/init.d/system
/etc/preinit
/lib/preinit/01_do_brcm6xxx_tch.sh
/etc/profile.d/bootloader-banner.sh
/etc/localtime
/lib/upgrade/platform.sh
/etc/shells
/etc/boards/VCNT-3/wlan/parse_pa_trimming.sh
/lib/functions/caldata.sh
/lib/preinit/70_initramfs_test
/etc/rc.button/reboot
/etc/device_info
/etc/boards/VCNT-3/wlan/brcm_dfs_cfg
/etc/protocols
/etc/uci-defaults/12_network-generate-ula
/lib/preinit/10_indicate_preinit
/usr/libexec/login.sh
/etc/boards/VCNT-D/wlan/bcmcmn_nvramvars.bin
/etc/rdpa_init.sh
/etc/fstab
/usr/lib/os-release
/etc/sysctl.conf
/sbin/sysupgrade
/lib/preinit/99_10_run_init
/sbin/led.sh
/etc/iproute2/rt_protos
/lib/preinit/50_indicate_regular_preinit
/etc/banner
/etc/shadow
/etc/rc.common
/bin/restricted_shell
/etc/uci-defaults/10_migrate-shadow
/lib/preinit/02_sysinfo
/etc/wlan
/lib/mount_root/10_remove_due_to_upgrade
/etc/board.d/99-default_network
/etc/os-release
/etc/board.d/02_network
/etc/boards/VCNT-3/wlan/bcmcmn_nvramvars.bin
/lib/preinit/80_mount_root
/etc/rc.local
/lib/functions/uci-defaults.sh
/etc/init.d/led
/etc/ethers
/etc/init.d/gpio_switch
/lib/functions/network.sh
/etc/init.d/sysctl
/sbin/hotplug-call
/usr/libexec/validate_firmware_image
/etc/boards/VCNT-3/wlan/brcm_country_map_2G
/lib/upgrade/fwtool.sh
/etc/services
/lib/functions/system.sh
/bin/board_detect
/lib/upgrade/common.sh
/etc/opkg/distfeeds.conf
/etc/boards/VCNT-3/wlan/VCNT-3_nvram.nvm
/etc/rc.button/power
/etc/sysupgrade.conf
/etc/boards/VCNT-3/wlan/brcm_country_map_5G
/etc/rc.button/rfkill
/lib/upgrade/keep.d/base-files
/lib/upgrade/do_stage2
/etc/group
/etc/rc.button/reset
/lib/upgrade/stage2
/etc/uci-defaults/13_fix_group_user
/lib/upgrade/rbi_vrss.lua
/sbin/firstboot
/etc/boards/VCNT-3/wlan/nvram2rip.sh
/etc/hotplug.d/net/00-sysctl
/etc/hosts
/lib/functions.sh
/etc/iproute2/rt_tables
/lib/functions/leds.sh
/etc/init.d/bcm-usb-support
/rom/note
/etc/inittab
/etc/openwrt_release
/sbin/wifi
/sbin/pkg_check
/etc/profile
/etc/boards/VCNT-D/wlan/brcm_country_map_2G
/etc/init.d/sysfixtime
/etc/TZ
/etc/init.d/umount
/lib/functions/service.sh
/lib/upgrade/keep.d/base-files-essential
/lib/functions/migrations.sh
/etc/init.d/bcm-device-nodes
/etc/sysctl.d/10-default.conf
/etc/passwd
/lib/preinit/02_default_set_state
/etc/iproute2/ematch_map
/etc/init.d/sensors
/etc/mtab
/lib/brcm6xxx_tch.sh
/etc/init.d/boot
/etc/boards/VCNT-D/wlan/brcm_country_map_5G
/etc/boards/VCNT-D/wlan/bcm943684_nvram.nvm
/etc/init.d/brcm_kthrd_prios
/usr/bin/odm_test_init.sh
/lib/functions/preinit.sh
/bin/ipcalc.sh
/etc/diag.sh
/lib/functions/functions-tch.sh
