# usr/lib/opkg/info/base-files.list
/etc/boards/VBNT-V/wlan/brcm_country_map_2G
/etc/boards/VCNT-X/wlan/brcm_dfs_cfg
/etc/boards/VBNT-6/wlan/bcm43602_map.bin
/lib/functions/syslog.sh
/etc/openwrt_version
/etc/hotplug.d/usb/12-tch_usb
/etc/banner.failsafe
/etc/resolv.conf
/etc/boards/VCNT-C/wlan/brcm_country_map_2G
/etc/boards/VCNT-A/wlan/brcm_dfs_cfg
/etc/boards/VBNT-6/wlan/bcm43664_nvramvars.bin
/etc/boards/VCNT-H/wlan/brcm_country_map_5G
/etc/init.d/done
/bin/config_generate
/var
/etc/init.d/bcm_modules
/etc/boards/VBNT-6/wlan/brcm_country_map_5G
/etc/opkg/keys/de98a2dd1d0f8a07
/etc/init.d/system
/etc/preinit
/etc/boards/VBNT-V/wlan/brcm_country_map_5G
/lib/upgrade/platform.sh
/etc/boards/VCNT-C/rdpa_init.sh
/etc/shells
/lib/preinit/10_sysinfo
/etc/boards/VCNT-C/wlan/brcm_country_map_5G
/lib/preinit/70_initramfs_test
/etc/device_info
/etc/boards/VCNT-E/wlan/bcm43602_map.bin
/etc/protocols
/etc/boards/VBNT-6/wlan/brcm_dfs_cfg
/etc/uci-defaults/12_network-generate-ula
/etc/boards/VCNT-H/rdpa_init.sh
/etc/boards/VCNT-C/wlan/brcm_dfs_cfg
/lib/preinit/10_indicate_preinit
/etc/boards/VBNT-V/wlan/brcm_dfs_cfg
/etc/boards/VCNT-E/wlan/bcm43664_nvramvars.bin
/etc/boards/VBNT-6/wlan/bcm43664_map.bin
/etc/boards/VCNT-H/wlan/bcm43602_map.bin
/lib/functions/uci-defaults-new.sh
/etc/boards/VBNT-6/init.d/bcm_modules
/etc/rdpa_init.sh
/etc/config/system
/usr/bin/dhcp6c-state
/etc/fstab
/etc/sysctl.conf
/etc/boards/VCNT-E/wlan/brcm_dfs_cfg
/sbin/sysupgrade
/etc/boards/VBNT-Y/wlan/bcm43664_map.bin
/etc/boards/VBNT-6/rdpa_init.sh
/lib/preinit/99_10_run_init
/etc/boards/VBNT-9/wlan/brcm_country_map_2G
/etc/opkg/keys/53bad1233d4c98c5
/sbin/led.sh
/etc/iproute2/rt_protos
/lib/preinit/50_indicate_regular_preinit
/etc/banner
/etc/boards/VBNT-V/wlan/bcm4366_map.bin
/etc/boards/VCNT-E/wlan/bcm43664_map.bin
/etc/boards/VBNT-Y/wlan/brcm_country_map_2G
/etc/shadow
/etc/boards/VCNT-A/wlan/bcm4363_nvramvars.bin
/etc/rc.common
/bin/restricted_shell
/etc/boards/VBNT-9/wlan/bcm43217_map.bin
/etc/uci-defaults/10_migrate-shadow
/etc/wlan
/etc/jffs2_options
/init
/etc/boards/VCNT-H/wlan/bcm43664_map.bin
/etc/boards/VBNT-9/wlan/brcm_country_map_5G
/etc/boards/VBNT-O/wlan/bcm43664_map.bin
/etc/boards/VCNT-X/wlan/bcm43664_map.bin
/etc/boards/VBNT-J/wlan/brcm_dfs_cfg
/lib/preinit/80_mount_root
/etc/boards/VCNT-C/wlan/bcm43664_nvramvars.bin
/etc/boards/VBNT-Y/wlan/brcm_country_map_5G
/etc/rc.local
/lib/functions/uci-defaults.sh
/etc/boards/VBNT-J/wlan/bcm4366_map.bin
/etc/boards/VCNT-X/wlan/brcm_country_map_2G
/etc/init.d/led
/etc/boards/VCNT-A/wlan/brcm_country_map_2G
/lib/functions/network.sh
/etc/init.d/sysctl
/sbin/hotplug-call
/etc/boards/VBNT-O/wlan/bcm43217_map.bin
/etc/boards/VBNT-O/wlan/brcm_country_map_2G
/etc/rc.button/failsafe
/etc/boards/VBNT-Y/wlan/bcm4363_map.bin
/etc/services
/etc/boards/VCNT-X/wlan/brcm_country_map_5G
/etc/boards/VCNT-A/wlan/bcm4363_map.bin
/lib/functions/system.sh
/etc/boards/VCNT-A/wlan/brcm_country_map_5G
/bin/board_detect
/etc/boards/VCNT-A/wlan/bcm43664_map.bin
/lib/upgrade/common.sh
/etc/boards/VCNT-E/wlan/brcm_country_map_2G
/etc/boards/VBNT-J/wlan/brcm_country_map_2G
/etc/sysupgrade.conf
/etc/rc.button/power
/etc/boards/VBNT-O/wlan/brcm_country_map_5G
/etc/rc.button/rfkill
/lib/upgrade/keep.d/base-files
/etc/boards/VBNT-Y/wlan/bcm4363_nvramvars.bin
/etc/boards/VCNT-H/init.d/bcm_modules
/etc/group
/etc/rc.button/reset
/etc/boards/VCNT-X/wlan/bcm4363_map.bin
/etc/boards/VCNT-A/wlan/bcm43664_nvramvars.bin
/etc/boards/VCNT-C/init.d/bcm_modules
/sbin/firstboot
/lib/upgrade/rbi_vrss.lua
/etc/boards/VCNT-E/wlan/brcm_country_map_5G
/etc/boards/VBNT-Y/wlan/bcm43664_nvramvars.bin
/etc/boards/VBNT-7/rdpa_init.sh
/etc/boards/VBNT-J/wlan/brcm_country_map_5G
/bin/login.sh
/etc/boards/VBNT-7/wlan/bcm43664_wl1_map.bin
/etc/hotplug.d/net/00-sysctl
/etc/boards/VBNT-W/wlan/bcm43602_map.bin
/etc/hosts
/etc/boards/VBNT-W/rdpa_init.sh
/lib/functions.sh
/etc/boards/VBNT-7/wlan/brcm_country_map_2G
/etc/uci-defaults/11_migrate-sysctl
/etc/config/network
/lib/functions/leds.sh
/etc/boards/VCNT-C/wlan/bcm43602_map.bin
/etc/boards/VBNT-9/wlan/bcm4360_map.bin
/etc/boards/VCNT-X/wlan/bcm4363_nvramvars.bin
/etc/boards/VBNT-W/wlan/brcm_country_map_2G
/rom/note
/etc/init.d/bcm-usb-support
/etc/inittab
/etc/openwrt_release
/sbin/wifi
/etc/boards/VCNT-H/wlan/bcm43664_nvramvars.bin
/etc/profile
/etc/boards/VBNT-7/init.d/bcm_modules
/etc/boards/VCNT-X/wlan/bcm43664_nvramvars.bin
/etc/boards/VBNT-V/wlan/bcm4360_map.bin
/etc/boards/VBNT-7/wlan/brcm_country_map_5G
/etc/init.d/sysfixtime
/etc/TZ
/etc/boards/VBNT-W/init.d/bcm_modules
/etc/init.d/umount
/lib/functions/service.sh
/etc/boards/VCNT-E/init.d/bcm_modules
/etc/boards/VBNT-W/wlan/brcm_country_map_5G
/lib/upgrade/keep.d/base-files-essential
/etc/boards/VBNT-Y/wlan/brcm_dfs_cfg
/etc/init.d/bcm-device-nodes
/etc/passwd
/etc/boards/VBNT-7/wlan/bcm43664_map.bin
/lib/preinit/02_default_set_state
/etc/init.d/sensors
/etc/mtab
/etc/init.d/boot
/etc/boards/VBNT-W/wlan/bcm43664_map.bin
/etc/boards/VBNT-7/wlan/bcm43664_nvramvars.bin
/etc/init.d/brcm_kthrd_prios
/usr/bin/odm_test_init.sh
/lib/functions/preinit.sh
/etc/boards/VCNT-E/rdpa_init.sh
/etc/boards/VBNT-J/wlan/bcm4360_map.bin
/bin/ipcalc.sh
/etc/boards/VCNT-H/wlan/brcm_country_map_2G
/etc/boards/VBNT-W/wlan/bcm43664_nvramvars.bin
/etc/diag.sh
/etc/boards/VCNT-C/wlan/bcm43664_map.bin
/etc/boards/VBNT-6/wlan/brcm_country_map_2G
/lib/functions/functions-tch.sh
/etc/boards/VCNT-H/wlan/brcm_dfs_cfg
