# usr/lib/opkg/info/base-files.list
/lib/functions/syslog.sh
/etc/openwrt_version
/etc/hotplug.d/usb/12-tch_usb
/etc/banner.failsafe
/etc/resolv.conf
/etc/boards/VBNT-H/wlan/brcm_country_map_2G
/etc/init.d/done
/bin/config_generate
/var
/etc/init.d/bcm_modules
/etc/opkg/keys/de98a2dd1d0f8a07
/etc/init.d/system
/etc/preinit
/lib/upgrade/platform.sh
/etc/shells
/lib/preinit/10_sysinfo
/lib/preinit/70_initramfs_test
/etc/device_info
/etc/protocols
/etc/uci-defaults/12_network-generate-ula
/lib/preinit/10_indicate_preinit
/etc/boards/VANT-W/wlan/bcm43217_map.bin
/lib/functions/uci-defaults-new.sh
/etc/rdpa_init.sh
/etc/config/system
/etc/boards/VBNT-S/wlan/bcm43217_map.bin
/usr/bin/dhcp6c-state
/etc/fstab
/etc/sysctl.conf
/sbin/sysupgrade
/lib/preinit/99_10_run_init
/etc/boards/VBNT-F/wlan/bcm43217_map.bin
/etc/opkg/keys/53bad1233d4c98c5
/sbin/led.sh
/lib/preinit/50_indicate_regular_preinit
/etc/iproute2/rt_protos
/etc/banner
/etc/shadow
/etc/rc.common
/bin/restricted_shell
/etc/uci-defaults/10_migrate-shadow
/etc/wlan
/etc/jffs2_options
/init
/etc/boards/VBNT-K/wlan/brcm_country_map_2G
/lib/preinit/80_mount_root
/etc/rc.local
/lib/functions/uci-defaults.sh
/etc/init.d/led
/etc/boards/VBNT-F/wlan/brcm_country_map_2G
/lib/functions/network.sh
/etc/init.d/sysctl
/sbin/hotplug-call
/etc/rc.button/failsafe
/etc/services
/lib/functions/system.sh
/bin/board_detect
/lib/upgrade/common.sh
/etc/sysupgrade.conf
/etc/rc.button/power
/etc/rc.button/rfkill
/lib/upgrade/keep.d/base-files
/etc/boards/VBNT-S/wlan/brcm_country_map_2G
/etc/group
/etc/rc.button/reset
/etc/boards/VBNT-H/wlan/bcm43217_map.bin
/lib/upgrade/rbi_vrss.lua
/sbin/firstboot
/bin/login.sh
/etc/hotplug.d/net/00-sysctl
/etc/hosts
/lib/functions.sh
/etc/boards/VBNT-K/wlan/bcm43217_map.bin
/etc/uci-defaults/11_migrate-sysctl
/etc/config/network
/lib/functions/leds.sh
/rom/note
/etc/init.d/bcm-usb-support
/etc/inittab
/etc/openwrt_release
/sbin/wifi
/etc/profile
/etc/init.d/sysfixtime
/etc/TZ
/lib/functions/service.sh
/etc/init.d/umount
/etc/boards/VANT-W/wlan/brcm_country_map_2G
/lib/upgrade/keep.d/base-files-essential
/etc/init.d/bcm-device-nodes
/etc/passwd
/lib/preinit/02_default_set_state
/etc/mtab
/etc/init.d/boot
/usr/bin/odm_test_init.sh
/etc/init.d/brcm_kthrd_prios
/lib/functions/preinit.sh
/bin/ipcalc.sh
/etc/diag.sh
/lib/functions/functions-tch.sh
