# usr/lib/opkg/info/base-files.list
/lib/preinit/31_create_device_nodes
/lib/functions/syslog.sh
/etc/openwrt_version
/etc/hotplug.d/usb/12-tch_usb
/etc/resolv.conf
/etc/wlan/bcm43217_map.bin
/etc/init.d/done
/var
/etc/init.d/bcm_modules
/lib/firstboot/99_10_no_fo_cleanup
/etc/hotplug2-init.rules
/etc/preinit
/lib/firstboot/50_pivot
/sbin/mount_root
/lib/firstboot/30_reset_copy_rom
/lib/upgrade/platform.sh
/etc/hotplug2-common.rules
/etc/shells
/lib/preinit/70_initramfs_test
/lib/preinit/50_choose_console
/lib/firstboot/10_no_fo_clear_overlay
/etc/protocols
/lib/preinit/20_check_jffs2_ready
/lib/preinit/10_indicate_preinit
/lib/firstboot/99_10_with_fo_cleanup
/etc/rdpa_init.sh
/etc/config/system
/usr/bin/dhcp6c-state
/etc/fstab
/lib/preinit/30_device_fs_daemons
/lib/preinit/90_mount_no_jffs2
/etc/sysctl.conf
/sbin/sysupgrade
/lib/preinit/70_pivot_jffs2_root
/lib/preinit/99_10_run_init
/lib/firstboot/20_has_mini_fo
/lib/firstboot/10_reset_has_mini_fo
/lib/firstboot/10_determine_parts
/sbin/led.sh
/lib/preinit/50_indicate_regular_preinit
/lib/firstboot/40_no_fo_copy_ramoverlay
/etc/banner
/lib/preinit/10_essential_fs
/etc/shadow
/etc/rc.common
/lib/preinit/40_init_shm
/etc/jffs2_options
/init
/lib/preinit/40_mount_jffs2
/lib/preinit/20_device_fs_mount
/lib/preinit/80_mount_root
/etc/init.d/sysntpd
/etc/rc.local
/lib/functions/uci-defaults.sh
/etc/init.d/led
/lib/functions/network.sh
/etc/init.d/sysctl
/lib/firstboot/30_is_rootfs_mounted
/sbin/hotplug-call
/etc/init.d/usb
/lib/modules/3.4.11-rt19
/lib/preinit/40_mount_devpts
/lib/preinit/90_init_console
/etc/services
/lib/preinit/99_10_mount_no_mtd
/lib/preinit/10_check_for_mtd
/lib/upgrade/common.sh
/etc/sysupgrade.conf
/lib/firstboot/05_firstboot_skip
/lib/upgrade/keep.d/base-files
/etc/uci-defaults/migrate-sysctl
/etc/wlan/brcm_country_map_2G
/etc/group
/etc/functions.sh
/lib/upgrade/rbi_vrss.lua
/sbin/firstboot
/bin/login.sh
/etc/uci-defaults/migrate-shadow
/etc/hosts
/etc/hotplug.d/usb/10-usb
/lib/functions.sh
/lib/preinit/60_init_hotplug
/usr/bin/sysupgrade-safe
/etc/init.d/rcS
/rom/note
/etc/inittab
/etc/openwrt_release
/lib/preinit/05_mount_skip
/sbin/wifi
/etc/profile
/etc/TZ
/lib/firstboot/20_no_fo_mount_jffs
/lib/functions/service.sh
/etc/init.d/umount
/lib/preinit/90_restore_config
/lib/upgrade/keep.d/base-files-essential
/lib/firstboot/30_no_fo_pivot
/etc/passwd
/etc/hotplug.d/ieee1394/10-ieee1394
/lib/preinit/02_default_set_state
/etc/mtab
/lib/firstboot/40_copy_ramoverlay
/etc/init.d/boot
/lib/functions/boot.sh
/usr/bin/odm_test_init.sh
/etc/init.d/brcm_kthrd_prios
/bin/ipcalc.sh
/etc/diag.sh
/lib/functions/functions-tch.sh
/lib/firstboot/20_reset_clear_jffs
/lib/preinit/41_merge_overlay_hooks
/etc/uci-defaults/network
