# usr/lib/opkg/info/base-files.list
/lib/functions/syslog.sh
/etc/openwrt_version
/etc/boards/GBNT-2/wlan/bcm43217_nvramvars.bin
/etc/banner.failsafe
/etc/cms_entity_info.d/eid_bcm_omcipmd.txt
/etc/resolv.conf
/etc/boards/GBNT-Z/wlan/brcm_country_map_2G
/etc/cms_entity_info.d/eid_bcm_voip.txt
/etc/init.d/done
/bin/config_generate
/var
/etc/init.d/bcm_modules
/etc/boards/GBNT-2/wlan/bcm43664_nvramvars.bin
/etc/boards/GBNT-Z/wlan/bcm43217_nvramvars.bin
/etc/opkg/keys/de98a2dd1d0f8a07
/etc/cms_entity_info.d/symbol_table.txt
/etc/init.d/system
/etc/preinit
/etc/boards/GBNT-Z/wlan/bcm43664_nvramvars.bin
/etc/boards/GBNT-5/wlan/brcm_country_map_2G
/lib/upgrade/platform.sh
/etc/hotplug2-common.rules
/etc/shells
/usr/bin/leds_onoff.sh
/lib/preinit/10_sysinfo
/lib/preinit/70_initramfs_test
/etc/init.d/lilac
/etc/device_info
/etc/boards/GBNT-Z/wlan/brcm_country_map_5G
/etc/protocols
/etc/uci-defaults/12_network-generate-ula
/opt/scripts/setmem
/lib/preinit/10_indicate_preinit
/etc/cms_entity_info.d/eid_bcm_base.txt
/etc/boards/GBNT-5/wlan/brcm_country_map_5G
/usr/lib/gpon/gpon_cli.sh
/lib/functions/uci-defaults-new.sh
/etc/rdpa_init.sh
/etc/config/system
/opt/scripts/sb
/opt/scripts/sh
/usr/bin/dhcp6c-state
/usr/bin/gpon_mtu
/etc/fstab
/etc/sysctl.conf
/sbin/sysupgrade
/opt/scripts/sw
/lib/preinit/70_pivot_jffs2_root
/lib/preinit/99_10_run_init
/opt/scripts/spitr
/usr/bin/prbs.sh
/etc/opkg/keys/53bad1233d4c98c5
/etc/boards/GBNT-Z/wlan/bcm43664_map.bin
/sbin/led.sh
/etc/iproute2/rt_protos
/lib/preinit/50_indicate_regular_preinit
/etc/banner
/etc/shadow
/etc/rc.common
/bin/restricted_shell
/etc/uci-defaults/10_migrate-shadow
/etc/wlan
/opt/scripts/dumpmem
/init
/etc/rdpa_common_filter_init.sh
/usr/bin/gpon_cli
/etc/boards/GBNT-Z/wlan/bcm43217_map.bin
/usr/bin/hotspot-ssid_activation
/data
/lib/preinit/80_mount_root
/etc/rc.local
/lib/functions/uci-defaults.sh
/etc/rdpa_gpon_init.sh
/etc/init.d/led
/lib/functions/network.sh
/etc/init.d/sysctl
/sbin/hotplug-call
/etc/hotplug.d/platform/00-sfp
/etc/boards/GBNT-6/wlan/bcm43664_map.bin
/etc/cms_entity_info.d/eid_bcm_omcid.txt
/etc/rc.button/failsafe
/etc/ledfw/stateMachines.lua
/etc/services
/etc/boards/GBNT-6/wlan/bcm43664_nvramvars.bin
/lib/functions/system.sh
/bin/board_detect
/lib/upgrade/common.sh
/etc/boards/GBNT-6/wlan/bcm4360_map.bin
/etc/rc.button/power
/etc/sysupgrade.conf
/etc/rc.button/rfkill
/lib/upgrade/keep.d/base-files
/usr/bin/lanleds_onoff.sh
/etc/group
/etc/rc.button/reset
/etc/boards/GBNT-5/wlan/bcm43664_nvramvars.bin
/lib/upgrade/rbi_vrss.lua
/sbin/firstboot
/etc/cms_entity_info.d/eid_bcm_mgmt.txt
/bin/login.sh
/etc/hotplug.d/net/00-sysctl
/etc/hosts
/lib/functions.sh
/etc/boards/GBNT-6/wlan/bcm4360_nvramvars.bin
/etc/uci-defaults/11_migrate-sysctl
/etc/config/network
/etc/boards/GBNT-2/wlan/brcm_country_map_2G
/lib/functions/leds.sh
/etc/boards/GBNT-2/wlan/bcm43664_map.bin
/etc/init.d/bcm-usb-support
/rom/note
/etc/inittab
/etc/openwrt_release
/opt/scripts/db
/opt/scripts/dh
/sbin/wifi
/opt/scripts/dw
/etc/profile
/usr/bin/gpon_intf_init
/opt/scripts/bcmlog
/etc/boards/GBNT-5/wlan/bcm43664_map.bin
/etc/init.d/sysfixtime
/etc/TZ
/etc/boards/GBNT-2/wlan/brcm_country_map_5G
/opt/scripts/fb
/etc/boards/GBNT-2/wlan/bcm43217_map.bin
/opt/scripts/fh
/etc/init.d/umount
/lib/functions/service.sh
/lib/upgrade/keep.d/base-files-essential
/etc/init.d/bcm-device-nodes
/opt/scripts/fw
/etc/passwd
/etc/boards/GBNT-6/wlan/brcm_country_map_2G
/lib/preinit/02_default_set_state
/etc/mtab
/etc/init.d/boot
/etc/boards/GBNT-5/wlan/bcm4360_nvramvars.bin
/etc/init.d/brcm_kthrd_prios
/usr/bin/odm_test_init.sh
/lib/functions/preinit.sh
/etc/boards/GBNT-5/wlan/bcm4360_map.bin
/etc/rdpa_common_init.sh
/bin/ipcalc.sh
/etc/diag.sh
/etc/cms_entity_info.d/eid_bcm_kthreads.txt
/lib/functions/functions-tch.sh
/etc/boards/GBNT-6/wlan/brcm_country_map_5G
