# usr/lib/opkg/info/base-files.list
/lib/functions/syslog.sh
/etc/boards/VANT-D/wlan/brcm_country_map_2G
/etc/openwrt_version
/etc/hotplug.d/usb/12-tch_usb
/etc/boards/VANT-6/wlan/bcm4360_map.bin
/etc/banner.failsafe
/etc/boards/VBNT-A/wlan/bcm4360_map.bin
/etc/resolv.conf
/etc/init.d/done
/bin/config_generate
/var
/etc/boards/VANT-R/wlan/brcm_country_map_5G
/etc/boards/VANT-6/wlan/brcm_country_map_2G
/etc/opkg/keys/de98a2dd1d0f8a07
/etc/boards/VANT-E/wlan/bcm6362_map.bin
/etc/init.d/system
/etc/preinit
/etc/boards/VANT-D/wlan/brcm_country_map_5G
/lib/upgrade/platform.sh
/etc/shells
/lib/preinit/10_sysinfo
/lib/preinit/70_initramfs_test
/etc/device_info
/etc/boards/VANT-6/wlan/brcm_country_map_5G
/etc/protocols
/etc/uci-defaults/12_network-generate-ula
/etc/boards/VBNT-L/wlan/brcm_country_map_2G
/lib/preinit/10_indicate_preinit
/etc/boards/VANT-6/wlan/bcm6362_map.bin
/etc/boards/VANT-D/wlan/bcm4360_map.bin
/lib/functions/uci-defaults-new.sh
/etc/boards/VBNT-A/wlan/bcm6362_map.bin
/etc/config/system
/etc/boards/VBNT-L/wlan/bcm4360_map.bin
/etc/boards/VBNT-L/wlan/brcm_country_map_5G
/usr/bin/dhcp6c-state
/etc/fstab
/etc/sysctl.conf
/sbin/sysupgrade
/lib/preinit/99_10_run_init
/etc/opkg/keys/53bad1233d4c98c5
/sbin/led.sh
/lib/preinit/50_indicate_regular_preinit
/etc/banner
/etc/boards/VANT-8/wlan/bcm4360_map.bin
/etc/shadow
/etc/rc.common
/etc/uci-defaults/10_migrate-shadow
/etc/wlan
/etc/jffs2_options
/init
/etc/boards/VANT-D/wlan/bcm6362_map.bin
/etc/boards/VANT-R/wlan/bcm4360_map.bin
/lib/preinit/80_mount_root
/etc/rc.local
/lib/functions/uci-defaults.sh
/etc/boards/VBNT-L/wlan/bcm6362_map.bin
/etc/init.d/led
/lib/functions/network.sh
/etc/init.d/sysctl
/sbin/hotplug-call
/etc/boards/VANT-4/wlan/brcm_country_map_2G
/lib/modules/3.4.11-rt19
/etc/boards/VANT-8/wlan/bcm6362_map.bin
/etc/boards/VANT-F/wlan/bcm4360_map.bin
/etc/rc.button/failsafe
/etc/services
/etc/boards/VBNT-A/wlan/brcm_country_map_2G
/lib/functions/system.sh
/etc/boards/VANT-F/wlan/brcm_country_map_2G
/bin/board_detect
/etc/boards/VANT-4/wlan/brcm_country_map_5G
/lib/upgrade/common.sh
/etc/boards/VANT-R/wlan/bcm6362_map.bin
/etc/sysupgrade.conf
/etc/rc.button/power
/etc/boards/VANT-8/wlan/brcm_country_map_2G
/etc/rc.button/rfkill
/lib/upgrade/keep.d/base-files
/etc/group
/etc/boards/VBNT-A/wlan/brcm_country_map_5G
/etc/rc.button/reset
/etc/boards/VANT-F/wlan/brcm_country_map_5G
/etc/boards/VANT-7/wlan/bcm4360_map.bin
/lib/upgrade/rbi_vrss.lua
/sbin/firstboot
/etc/boards/VANT-4/wlan/bcm4360_map.bin
/bin/login.sh
/etc/boards/VANT-8/wlan/brcm_country_map_5G
/etc/hotplug.d/net/00-sysctl
/etc/boards/VANT-F/wlan/bcm6362_map.bin
/etc/hosts
/lib/functions.sh
/etc/uci-defaults/11_migrate-sysctl
/lib/functions/leds.sh
/usr/bin/sysupgrade-safe
/rom/note
/etc/inittab
/etc/openwrt_release
/etc/boards/VANT-E/wlan/brcm_country_map_2G
/sbin/wifi
/etc/profile
/etc/boards/VANT-7/wlan/brcm_country_map_2G
/etc/init.d/sysfixtime
/etc/TZ
/lib/functions/service.sh
/etc/init.d/umount
/etc/boards/VANT-7/wlan/bcm6362_map.bin
/lib/upgrade/keep.d/base-files-essential
/etc/boards/VANT-E/wlan/bcm4360_map.bin
/etc/init.d/bcm-device-nodes
/etc/passwd
/etc/boards/VANT-E/wlan/brcm_country_map_5G
/lib/preinit/02_default_set_state
/etc/boards/VANT-4/wlan/bcm6362_map.bin
/etc/mtab
/etc/init.d/boot
/usr/bin/odm_test_init.sh
/etc/init.d/brcm_kthrd_prios
/etc/boards/VANT-7/wlan/brcm_country_map_5G
/lib/functions/preinit.sh
/bin/ipcalc.sh
/etc/diag.sh
/etc/boards/VANT-R/wlan/brcm_country_map_2G
/sbin/core-handler
/lib/functions/functions-tch.sh
/etc/uci-defaults/network
