# usr/lib/opkg/info/bcm-spdsvc-client.list
/usr/bin/bcm_spdsvc
/etc/init.d/bcm_spdsvc
