# usr/lib/opkg/info/bcm_userspace.list
/usr/bin/swmdk
/usr/lib/libbcm_boardctl.so
/etc/init.d/power
/usr/bin/iqctl
/usr/bin/tmctl
/usr/lib/libtmctl.so
/usr/lib/libspdsvc.so
/usr/bin/xtmctl
/usr/lib/libgen_util.so
/etc/init.d/license_init.sh
/usr/lib/libfcctl.so
/usr/lib/libvlanctl.so
/usr/bin/bpmctl
/usr/bin/tmsctl
/etc/init.d/phypower
/bin/runner
/usr/bin/dsldiagd
/lib/network/vopi_intfs_operation.sh
/usr/lib/libatmctl.so
/usr/lib/libxdslctl.so
/etc/init.d/spu
/usr/lib/libethctl.so
/usr/bin/vlanctl
/usr/lib/libvoip.so
/usr/bin/wlctl
/usr/bin/fcctl
/usr/lib/libspuctl.so
/usr/lib/lua/setuptm.lua
/usr/lib/libbdmf.so
/usr/lib/libiqctl.so
/etc/init.d/xdsl
/usr/lib/libbcmmcast.so
/usr/bin/mdkshell
/etc/init.d/wireless
/bin/bdmf_shell
/usr/bin/wl
/usr/bin/dhd
/usr/bin/xdslctl
/usr/lib/lua/tmhelper.lua
/etc/rc.d/S71license_init
/usr/bin/nvram
/usr/lib/libnvram.so
/usr/lib/libwlcsm.so
/bin/bs
/etc/init.d/swmdk
/etc/hotplug.d/dsl/10_trigger_xtm
/usr/lib/libbcm_flashutil.so
/usr/lib/libwlctl.so
/etc/init.d/vopi_intfs
/usr/bin/spuctl
/usr/bin/bcmswconfig
/usr/bin/dhdctl
/usr/lib/libpwrctl.so
/etc/licenses/license.txt
/usr/bin/ethctl
/lib/network/switch.sh
/etc/init.d/ethernet
/etc/hotplug.d/net/10_ethoam
/usr/lib/libethswctl.so
/usr/bin/pwrctl
/etc/init.d/iqos
/usr/lib/libbcm_util.so
/etc/init.d/xtm
/usr/lib/libsys_util.so
/etc/init.d/ethoam
/usr/bin/ethswctl
/usr/lib/librdpactl.so
/usr/lib/libbcm_spdt.so
/usr/bin/setuprunner
/etc/broadcom_release
/lib/functions/phypower.sh
/etc/init.d/softswitch
