# usr/lib/opkg/info/bcm_userspace.list
/usr/bin/swmdk
/usr/lib/libbcm_boardctl.so
/etc/init.d/power
/usr/bin/iqctl
/usr/bin/tmctl
/usr/lib/libtmctl.so
/usr/lib/libgen_util.so
/usr/lib/libfcctl.so
/usr/lib/libvlanctl.so
/usr/bin/bpmctl
/usr/bin/tmsctl
/etc/init.d/phypower
/lib/network/vopi_intfs_operation.sh
/usr/lib/libethctl.so
/usr/bin/vlanctl
/usr/bin/wlctl
/usr/bin/fcctl
/usr/lib/libarcher.so
/usr/bin/setuparcher
/usr/lib/lua/setuptm.lua
/usr/lib/libiqctl.so
/usr/lib/libbcmmcast.so
/usr/bin/mdkshell
/etc/init.d/wireless
/usr/bin/wl
/usr/bin/dhd
/usr/lib/lua/tmhelper.lua
/usr/bin/nvram
/usr/lib/libnvram.so
/usr/lib/libwlcsm.so
/etc/init.d/swmdk
/usr/lib/libbcm_flashutil.so
/etc/init.d/vopi_intfs
/usr/bin/bcmswconfig
/usr/bin/dhdctl
/usr/lib/libpwrctl.so
/usr/bin/ethctl
/lib/network/switch.sh
/etc/init.d/ethernet
/etc/hotplug.d/net/10_ethoam
/usr/lib/libethswctl.so
/usr/bin/archerctl
/usr/bin/pwrctl
/etc/init.d/iqos
/usr/lib/libbcm_util.so
/usr/lib/libsys_util.so
/etc/init.d/ethoam
/usr/bin/ethswctl
/etc/broadcom_release
/lib/functions/phypower.sh
/etc/init.d/softswitch
