# usr/lib/opkg/info/bcm_userspace-tch.list
/etc/init.d/phypower
/lib/network/switch.sh
/etc/init.d/ethernet
/lib/functions/phypower.sh
/etc/init.d/softswitch
