# usr/lib/opkg/info/bcm_userspace_pon-tch.list
/etc/init.d/bcm-usr
/usr/bin/gpon_tm_reload
/etc/init.d/gponl3
/lib/network/gpon.sh
/etc/init.d/bcm-wanconf
/etc/init.d/gpon-uni
/etc/init.d/gpon-ethwan
