# usr/lib/opkg/info/bcmubusbridge.list
/etc/init.d/bcmubusbridge
/usr/bin/bcmubusbridge
