# usr/lib/opkg/info/bulkdata.list
/etc/hotplug.d/ntp/20-bulkdata
/etc/init.d/bulkdata
/usr/bin/bulkdata
