# usr/lib/opkg/info/cgevent.list
/usr/bin/cgevent
/etc/init.d/cgevent
/lib/upgrade/keep.d/cgevent
/lib/cgevent/cups
