# usr/lib/opkg/info/conf-lxc.list
/srv/lxc/lxc_ee/rootfs/etc/init.d/zend
/usr/share/lxc/templates/lxc-homewarebox
/srv/lxc/lxc_ee/rootfs/etc/rc.d/S05lxc_start
/usr/share/lxc/templates/lxc-dropbear
/usr/share/lxc/templates/lxc-lcm
/srv/lxc/lxc_ee/rootfs/etc/init.d/lxc_start
/srv/lxc/lxc_ee/rootfs/etc/rc.d/S99zend
/etc/custo/config/lxc
/srv/lxc/lxc_ee/rootfs/etc/rc.d/K01zend
