# usr/lib/opkg/info/conf-system.list
/etc/custo/config/system
/usr/bin/memorystats.sh
/etc/uci-defaults/tch_0090-system
/etc/init.d/memorystats
/usr/lib/lua/memorystats-smap.lua
/.init_enable_core
/etc/memorystats.d/memorystats.crontab
