# usr/lib/opkg/info/conf-lxc.list
/usr/share/acl.d/ubusd_acl.json
/usr/share/lxc/include.conf
/usr/share/lxc/templates/lxc-homewarebox
/srv/lxc/lxc_ee/rootfs/etc/rc.d/S05lxc_start
/usr/share/lxc/templates/lxc-dropbear
/usr/share/lxc/templates/lxc-lcm
/srv/lxc/lxc_ee/rootfs/etc/init.d/lxc_start
/etc/custo/config/lxc
