# usr/lib/opkg/info/cujo-rabid-oem.list
/opt/cujo/lib/rabid-functions.sh
/opt/cujo/etc/tls-cujo-cloud.pem
/opt/cujo/share/lua/cujo/config/parameters.lua
/etc/init.d/rabid
/opt/cujo/bin/rabid-sh-env
