# usr/lib/opkg/info/cujo-luapermission.list
/etc/init.d/cujo-user
/opt/cujo/lib/lua/cujo/permission.so
