# usr/lib/opkg/info/cwmpd.list
/usr/lib/cwmpd/transfers/switchover.sh
/usr/lib/cwmpd/transfers/rollback_common.sh
/usr/lib/cwmpd/libuciplugin.so
/usr/lib/cwmpd/transfers/upgrade.sh
/usr/lib/cwmpd/transfers/handle_image.sh
/usr/lib/cwmpd/transfers/common_functions.sh
/usr/lib/cwmpd/transfers/target/gateway/start
/usr/lib/cwmpd/transfers/target/mobiled/cleanup
/usr/lib/cwmpd/transfers/vendorconfig.sh
/usr/bin/cwmpd
/usr/lib/cwmpd/transfers/rollback.sh
/usr/lib/cwmpd/transfers/target/mobiled/get_error
/usr/lib/cwmpd/transfers/target/gateway/get_error
/usr/lib/cwmpd/transfers/rollback_check.sh
/usr/lib/cwmpd/transfers/exec_config.lua
/usr/bin/cwmpevents
/usr/lib/cwmpd/transfers/vendorlog.sh
/etc/init.d/cwmpd
/usr/lib/cwmpd/transfers/config.sh
/lib/functions/firewall-cwmpd.sh
/usr/lib/cwmpd/transfers/target/gateway/cleanup
/usr/lib/libcwmpTransformerBackend.so
/usr/lib/cwmpd/transfers/cwmpd_parameter_conversion.sh
/usr/lib/cwmpd/transfers/target/mobiled/start
/usr/lib/cwmpd/libtransferplugin.so
/usr/lib/cwmpd/transfers/checkOngoingServices.lua
