# usr/lib/opkg/info/datausage-monitor.list
/usr/bin/datausaged
/etc/hotplug.d/ntp/35-datausaged
/etc/init.d/datausaged
