# usr/lib/opkg/info/datausage-notifier.list
/usr/bin/datausage_notifier
/usr/lib/lua/datausage_notifier/intercept.lua
/usr/lib/lua/datausage_notifier/sms.lua
/etc/init.d/datausage_notifier
/usr/lib/lua/datausage_notifier/email.lua
