# usr/lib/opkg/info/dhcpsnooper-tch.list
/usr/bin/dhcpsnooper
/etc/uci-defaults/70-dhcpsnooper
/etc/init.d/dhcpsnooper
/lib/upgrade/keep.d/dhcpsnooper-tch
/lib/functions/firewall-dhcpsnooper.sh
