# usr/lib/opkg/info/dhcpopassthrud.list
/lib/upgrade/keep.d/dhcpopassthrud
/etc/init.d/dhcpopassthrud
/usr/bin/dhcpopassthrud
