# usr/lib/opkg/info/dhcprelay.list
/etc/init.d/dhcprelay
/etc/hotplug.d/iface/30-dhcprelay
