# usr/lib/opkg/info/dlnad.list
/usr/bin/dlnad
/usr/lib/libDmsIntegration.so
/usr/lib/libXbox.so
/usr/lib/libILibPluginCommunication.so
/usr/lib/libProtocolInfoParser.so
/usr/lib/libMediaServerAbstraction.so
/usr/lib/libDownloadController.so
/etc/init.d/dlnad
/usr/lib/libStringUtils.so
/usr/lib/libDmsLib.so
/usr/lib/libMbusClient.so
/usr/lib/libMvfs.so
/usr/lib/libFileSystemEnumerator.so
/usr/lib/libHttpFiles.so
/usr/lib/libVirtualDirectory.so
/usr/lib/libCdsObjects.so
/usr/lib/libFileAbstractionFiles.so
