# usr/lib/opkg/info/dnsfilterd-tch.list
/etc/init.d/dnsfilter
/usr/sbin/dnsfilterd
