# usr/lib/opkg/info/dnsmasq-dhcpv6.list
/etc/dnsmasq.conf
/etc/hotplug.d/iface/25-dnsmasq
/etc/config/dhcp
/lib/dnsmasq/dhcp-event.sh
/etc/init.d/dnsmasq
/usr/sbin/dnsmasq
