# usr/lib/opkg/info/dnsmasq-full.list
/etc/dnsmasq.conf
/etc/config/dhcp
/lib/dnsmasq/dhcp-event.sh
/etc/hotplug.d/ntp/25-dnsmasqsec
/etc/init.d/dnsmasq
/usr/sbin/dnsmasq
