# usr/lib/opkg/info/dropbear.list
/usr/bin/dbclient
/etc/dropbear/dropbear_rsa_host_key
/etc/init.d/dropbear
/usr/bin/dropbearkey
/etc/dropbear/dropbear_ecdsa_host_key
/etc/config/dropbear
/usr/sbin/dropbear
