# usr/lib/opkg/info/ethctl.list
/etc/init.d/phypower
/usr/bin/ethctl
/etc/init.d/ethernet
/lib/functions/phypower.sh
