# usr/lib/opkg/info/faultmgmt.list
/usr/bin/processerror
/etc/init.d/faultmgmt
/usr/lib/lua/faultmgmt_event_monitor.lua
/usr/bin/faultmgmt.lua
/etc/hotplug.d/ntp/35-unsync
