# usr/lib/opkg/info/firewall-ext-tch-pinholehelper.list
/etc/init.d/pinholehelper
/usr/bin/pinholehelper.lua
