# usr/lib/opkg/info/firewall-ext-tch-redirecthelper.list
/usr/bin/redirecthelper.lua
/etc/init.d/redirecthelper
