# usr/lib/opkg/info/firmware-upgrade-ext.list
/usr/sbin/firmware-upgrade-ext.sh
/etc/init.d/firmware-upgrade-ext
