# usr/lib/opkg/info/firstusedate.list
/sbin/firstusedate
/etc/init.d/dateboot
