# usr/lib/opkg/info/fseventd.list
/usr/lib/libevprod.so
/usr/lib/libevlisten.so
/etc/init.d/fseventd
/usr/bin/fseventd
