# usr/lib/opkg/info/gre-hotspotd.list
/usr/bin/gre-hotspotd.lua
/usr/bin/gre-hotspotd-reload
/etc/init.d/gre-hotspotd
