# usr/lib/opkg/info/hostmanager.list
/etc/init.d/hostmanager
/usr/bin/hostmanager.lua
/etc/hotplug.d/ntp/05-hostmanager
/etc/config/user_friendly_name
