# usr/lib/opkg/info/hotspotd.list
/etc/init.d/hotspotd
/usr/sbin/hotspotd
/usr/sbin/fonctl
