# usr/lib/opkg/info/igmpproxy-tch.list
/etc/init.d/igmpproxy
/usr/lib/igmpproxy/libuciubusplugin.so
/usr/bin/igmpproxy
/etc/config/igmpproxy
