# usr/lib/opkg/info/igmpproxy-tch.list
/etc/init.d/igmpproxy
/usr/lib/igmpproxy/libuciubusplugin.so
/etc/uci-defaults/70-igmpproxy-tch
/lib/functions/firewall-igmpproxy-tch.sh
/usr/bin/igmpproxy
/etc/config/igmpproxy
