# usr/lib/opkg/info/igmpproxy-tch.list
/etc/init.d/igmpproxy
/usr/lib/igmpproxy/libuciubusplugin.so
/lib/upgrade/keep.d/igmpproxy-tch
/usr/sbin/igmpproxy
