# usr/lib/opkg/info/intercept-tch.list
/etc/init.d/intercept
/usr/lib/intercept/functions.sh
/etc/hotplug.d/iface/25-intercept
/usr/lib/intercept/setup.sh
/usr/lib/intercept/firewall.sh
