# usr/lib/opkg/info/intercept-tch.list
/etc/init.d/intercept
/usr/lib/intercept/common.sh
/usr/bin/interceptd
/usr/lib/intercept/intercept.sh
/usr/lib/intercept/plugins/wan.sh
/usr/lib/intercept/firewall.sh
