# usr/lib/opkg/info/ipset-helper-tch.list
/etc/init.d/ipset
/usr/sbin/ipset-helper
/lib/upgrade/keep.d/ipset-helper-tch
