# usr/lib/opkg/info/kmod-platform-gpl.list
/etc/init.d/platform-gpl
/etc/modules.d/51-platform-gpl
/lib/modules/4.1.52/platform-gpl.ko
