# usr/lib/opkg/info/kmodparams.list
/lib/upgrade/keep.d/kmodparams
/etc/init.d/kmodparams
