# usr/lib/opkg/info/kmod-platform.list
/etc/init.d/prozone
/lib/modules/3.4.11/platform.ko
/etc/hotplug.d/iface/99-uploadcore
/usr/bin/boot_srr
/lib/functions/reboot_reason.sh
/etc/modules.d/99-platform
/etc/init.d/warmboot
/etc/syslog_fwd/syslog_fwd.awk
/usr/sbin/kernel-core-handler
/etc/config/kernel_crash_handler
/etc/syslog_fwd/start_fwd_process
/usr/bin/bootp
/etc/hotplug.d/iface/99-warmboot
/etc/hotplug.d/iface/99-syslog-fwd
/etc/init.d/syslog_fwd
/etc/syslog_fwd/syslogd_restart
