# usr/lib/opkg/info/kmod-platform-gpl.list
/etc/init.d/platform-gpl
/lib/modules/3.4.11/platform-gpl.ko
