# usr/lib/opkg/info/kmod-platform.list
/etc/init.d/prozone
/usr/bin/boot_srr
/lib/functions/reboot_reason.sh
/etc/hotplug.d/iface/99-kernelpanic
/etc/modules.d/52-platform
/sbin/kernel-panic-handler
/etc/init.d/warmboot
/etc/syslog_fwd/syslog_fwd.awk
/lib/modules/4.1.38/platform.ko
/etc/syslog_fwd/start_fwd_process
/usr/bin/bootp
/etc/hotplug.d/iface/99-warmboot
/etc/hotplug.d/iface/99-syslog-fwd
/etc/init.d/syslog_fwd
/etc/syslog_fwd/syslogd_restart
/sbin/core-handler
