# usr/lib/opkg/info/lcmd.list
/etc/init.d/lcmd
/usr/lib/lua/lcm/execenv/native/execmon.lua
/usr/lib/lua/lcm/execenv.lua
/usr/lib/lua/lcm/execenv/native/opkg.lua
/usr/lib/lua/lcm/execenv/native/ipk.lua
/usr/lib/lua/lcm/execenv/opkg.lua
/usr/lib/lua/lcm/execenv/native/lxc_pid.lua
/usr/lib/lua/lcm/execenv/environment.lua
/usr/lib/lua/lcm/package.lua
/usr/lib/lua/lcm/store/storage.lua
/usr/lib/lua/lcm/execenv/lxc_opkg.lua
/usr/lib/lua/lcm/ubus.lua
/usr/lib/lua/lcm/store/filepath.lua
/usr/lib/lua/lcm/execenv/opkg_common.lua
/etc/public.pem
/usr/lib/lua/lcm/db.lua
/usr/lib/lua/lcm/queue.lua
/usr/lib/lua/lcm/ee_action.lua
/usr/lib/lua/lcm/execenv/native/common.lua
/usr/lib/lua/lcm/execenv/ee_loader.lua
/usr/lib/lua/lcm/store/file_store.lua
/usr/lib/lua/lcm/core.lua
/usr/lib/lua/lcm/execenv/native/opkg_package.lua
/usr/lib/lua/lcm/state_machine.lua
/opt/usr/lib/opkg/status
/usr/sbin/ee_action.lua
/usr/bin/lcm
/usr/lib/lua/lcm/execenv/memory.lua
/usr/lib/lua/lcm/store/none_store.lua
/usr/lib/lua/lcm/errorcodes.lua
/usr/lib/lua/lcm/execute.lua
/usr/bin/lcmd
