# usr/lib/opkg/info/lcmd.list
/etc/init.d/lcmd
/usr/lib/lua/lcm/execenv.lua
/usr/share/transformer/mappings/uci/lcmd.map
/etc/config/lcmd
/usr/lib/lua/lcm/execenv/native/opkg.lua
/usr/lib/lua/lcm/execenv/native/ipk.lua
/usr/lib/lua/lcm/operations.lua
/usr/lib/lua/lcm/package.lua
/usr/lib/lua/lcm/ubus.lua
/etc/public.pem
/usr/lib/lua/lcm/db.lua
/usr/lib/lua/lcm/queue.lua
/usr/lib/lua/lcm/execenv/native/common.lua
/usr/lib/lua/lcm/execenv/native.lua
/usr/lib/lua/lcm/core.lua
/usr/share/transformer/commitapply/lcmd.ca
/usr/lib/lua/lcm/execenv/native/opkg_package.lua
/opt/usr/lib/opkg/status
/usr/bin/lcm
/usr/share/transformer/mappings/bbf/SoftwareModules.map
/usr/lib/lua/lcm/errorcodes.lua
/usr/lib/lua/lcm/execenv/base_system.lua
/usr/bin/lcmd
