# usr/lib/opkg/info/ledfw.list
/sbin/ledfw.lua
/usr/lib/lua/ledframework/statemachine.lua
/usr/lib/lua/ledframework/ledhelper.lua
/usr/lib/lua/ledframework/ledaction.lua
/etc/init.d/ledfw
/usr/lib/lua/ledframework/init.lua
/usr/lib/lua/ledframework/common.lua
/usr/lib/lua/ledframework/patternaction.lua
