# usr/lib/opkg/info/libcgroup-tch.list
/usr/sbin/cgclear
/usr/bin/cgsnapshot
/usr/bin/cgcreate
/usr/bin/cgget
/usr/bin/lscgroup
/usr/lib/libcgroup.so
/usr/bin/cgexec
/usr/lib/libcgroup.so.1.0.41
/usr/bin/lssubsys
/etc/init.d/cgred
/usr/lib/libcgroup.so.1
/usr/bin/cgclassify
/etc/cgrules.conf
/etc/init.d/cgconfig
/usr/sbin/cgconfigparser
/usr/bin/cgdelete
/etc/cgconfig.conf
/usr/bin/cgset
/usr/sbin/cgrulesengd
