# usr/lib/opkg/info/lotagent.list
/etc/init.d/lotagent
/usr/bin/LoTAgent.lua
