# usr/lib/opkg/info/lsyslogd.list
/etc/init.d/lsyslogd
/usr/lib/lua/luasyslog.lua
/usr/bin/luasyslog
