# usr/lib/opkg/info/lte-doctor.list
/usr/bin/lte-doctor-dump
/usr/bin/lte-doctor-logger
/etc/init.d/lte-doctor-logger
