# usr/lib/opkg/info/lte-doctor.list
/usr/lib/lua/ltedoctor/plugins/mobiled_serving_system_changes.lua
/usr/lib/lua/ltedoctor/plugins/mobiled_signal_quality_changes.lua
/usr/bin/lte-doctor-logger
/etc/init.d/lte-doctor-logger
