# usr/lib/opkg/info/lua-bcm.list
/etc/init.d/portmirror
/usr/bin/setportmirror
/usr/bin/setportmirrorlan2lan
/etc/config/portmirror
/usr/bin/setxdsloemparam
/usr/lib/lua/luabcm.so
