# usr/lib/opkg/info/lxc-tch.list
/etc/init.d/lxc
/usr/sbin/update_lxc_config
/sbin/core-handler-lxc
/etc/init.d/lxc.mount
