# usr/lib/opkg/info/lxc-tch.list
/etc/init.d/lxc
/usr/sbin/update_lxc_config
