# usr/lib/opkg/info/lxc-tch.list
/usr/share/lxc/includes/curl.config
/etc/init.d/lxc
/usr/sbin/update_lxc_config
/sbin/core-handler-lxc
/usr/lib/lxc/net/hook
/etc/init.d/lxc.mount
/usr/sbin/lxc_monitor.lua
