# usr/lib/opkg/info/mabr.list
/etc/uci-defaults/tch_5002_NanoCDN
/root/nanocdn.ipk
/root/temp_nanoCDN/nanocdn.conf
/etc/rc.d/S99zAcotel
/root/temp_nanoCDN/nanocdn
/etc/init.d/nanocdn
/etc/rc.d/S95done
