# usr/lib/opkg/info/mcsnooper-tch.list
/etc/init.d/mcsnooper
/lib/upgrade/keep.d/mcsnooper-tch
/usr/bin/mcsnooper
