# usr/lib/opkg/info/mini-tr64.list
/etc/hotplug.d/iface/50-minitr064d
/etc/init.d/minitr064d
/usr/sbin/computeHA1
/usr/sbin/minitr064d
