# usr/lib/opkg/info/minidlna.list
/usr/sbin/minidlnad
/usr/sbin/cgroup_wrapper.sh
/usr/sbin/minidlna_bg
/etc/init.d/minidlna
/etc/hotplug.d/usb/45-minidlna
/etc/config/minidlna
