# usr/lib/opkg/info/miniupnpd-tch.list
/usr/share/miniupnpd-tch/firewall.include
/etc/config/upnpd
/etc/hotplug.d/iface/50-miniupnpd-tch
/etc/uci-defaults/99-miniupnpd-tch
/etc/init.d/miniupnpd-tch
