# usr/lib/opkg/info/minidlna-tch.list
/usr/sbin/minidlnad
/usr/sbin/cgroup_wrapper.sh
/usr/sbin/minidlna_bg
/etc/init.d/minidlna
/etc/config/minidlna
/etc/init.d/minidlna-procd
/etc/hotplug.d/mount/10-minidlna
