# usr/lib/opkg/info/mldproxy-tch.list
/usr/lib/mldproxy/libuciubusplugin.so
/etc/init.d/mldproxy
/lib/upgrade/keep.d/mldproxy-tch
/usr/sbin/mldproxy
