# usr/lib/opkg/info/mldproxy-tch.list
/etc/uci-defaults/70-mldproxy-tch
/usr/lib/mldproxy/libuciubusplugin.so
/etc/init.d/mldproxy
/usr/bin/mldproxy
/etc/config/mldproxy
/lib/functions/firewall-mldproxy-tch.sh
