# usr/lib/opkg/info/mldproxy-tch.list
/usr/lib/mldproxy/libuciubusplugin.so
/etc/init.d/mldproxy
/usr/bin/mldproxy
/etc/config/mldproxy
