# usr/lib/opkg/info/mmbrcm.list
/etc/init.d/mmpbxd
/usr/lib/libmmbrcm.so
