# usr/lib/opkg/info/mmbrcm.list
/etc/init.d/mmpbxd
/etc/init.d/mmpbxd_lite
/usr/lib/libmmbrcm.so
