# usr/lib/opkg/info/mmpbxfwctl.list
/usr/bin/mmpbxfwctl
/etc/init.d/mmpbxfwctl
