# usr/lib/opkg/info/mmpbxbrcm-dect.list
/usr/lib/mmpbx/libmmpbxbrcmdect.so
/etc/init.d/mmpbxbrcmdect
