# usr/lib/opkg/info/mobiledongle.list
/etc/init.d/mobiledongle
/usr/lib/lua/mobiledongle/at_helper.lua
/usr/bin/mobiledongle.lua
/etc/gcom/at_cmd.gcom
/etc/hotplug.d/usb/20-modeswitch
/usr/lib/lua/mobiledongle/state_helper.lua
/etc/hotplug.d/usb/30-mobiledongle
/usr/bin/mobiledongle_bg.lua
/etc/hotplug.d/usb/15-update_ms_uci
/usr/lib/lua/mobiledongle/bg_async.lua
