# usr/lib/opkg/info/mobiled.list
/lib/netifd/mobiled.script
/usr/lib/lua/mobiled/scriptloader.lua
/usr/lib/lua/mobiled/Doxyfile
/usr/lib/lua/mobiled/platform.lua
/usr/lib/lua/mobiled/mobiled.lua
/usr/lib/lua/mobiled/init.lua
/usr/lib/lua/mobiled/scripthelpers.lua
/etc/hotplug.d/iface/40-mobiled
/usr/bin/ubus_wait_for_event
/etc/init.d/mobiled
/usr/lib/lua/mobiled/informer.lua
/usr/lib/lua/mobiled/ubus.lua
/usr/lib/lua/libplatform.lua
/usr/lib/lua/mobiled/plugin.lua
/etc/hotplug.d/usb/50-mobiled
/lib/netifd/proto/mobiled.sh
/usr/lib/lua/mobiled/sms.lua
/usr/lib/lua/mobiled/version.lua
/usr/bin/usbmode_config
/usr/lib/lua/mobiled/device.lua
/usr/lib/lua/mobiled/detector.lua
/lib/netifd/mobiled.lua
/etc/uci-defaults/tch_0080-usbmode
/usr/lib/lua/mobiled/state.lua
/usr/bin/mobiled
/usr/lib/lua/mobiled/config.lua
/usr/lib/lua/mobiled/statemachine.lua
/usr/lib/lua/mobiled/error.lua
/etc/chatscripts/mobiled.chat
/usr/lib/lua/mobiled/events.lua
