# usr/lib/opkg/info/mosquitto-ssl.list
/etc/init.d/mosquitto
/etc/mosquitto/mosquitto.conf
/usr/sbin/mosquitto
