# usr/lib/opkg/info/mountd.list
/etc/init.d/mountd
/etc/config/mountd
/sbin/mountd
