# usr/lib/opkg/info/mqttjson-services.list
/etc/init.d/mqttjson-services
/usr/lib/lua/mqttjson_services/mqtt_client.lua
/usr/sbin/mqttjson_services.lua
/usr/lib/lua/mqttjson_services/logger.lua
/usr/lib/lua/mqttjson_services/bind.lua
/usr/lib/lua/mqttjson_services/utils.lua
