# usr/lib/opkg/info/mud.list
/etc/init.d/mud
/lib/functions/dlnad.sh
/usr/bin/mud
/lib/functions/mud_config.sh
/etc/mud/mud_file_type_config.ini
