# usr/lib/opkg/info/mud.list
/etc/init.d/mud
/usr/bin/mud
/etc/mud/mud_file_type_config.ini
