# usr/lib/opkg/info/multiap_controller.list
/usr/bin/multiap_controller
/usr/bin/map_common_reset_credentials.sh
/usr/bin/map_reset_to_default.sh
/usr/bin/multiap-controller-uci-to-procd-env
/etc/init.d/multiap_controller
