# usr/lib/opkg/info/multiap_controller.list
/usr/bin/multiap_controller
/usr/bin/reset_ctrl_credentials.sh
/usr/bin/map_common_reset_credentials.sh
/usr/bin/map_reset_to_default.sh
/etc/init.d/multiap_controller
/usr/bin/prepare_certi_controller.sh
