# usr/lib/opkg/info/multiap_agent.list
/usr/bin/config_bss_agent.sh
/etc/init.d/multiap_agent
/usr/bin/prepare_certi_agent.sh
/usr/bin/multiap_agent
