# usr/lib/opkg/info/multiap_vendorextensions.list
/usr/lib/lua/vendorextensions/client_main.lua
/usr/lib/lua/vendorextensions/msg.lua
/usr/lib/lua/vendorextensions/actionHandler.lua
/usr/lib/lua/vendorextensions/msg_decode.lua
/usr/lib/lua/vendorextensions/uciconfig.lua
/usr/lib/lua/vendorextensions/init.lua
/usr/bin/multiap_vendorextensions
/usr/lib/lua/vendorextensions/cron_handler.lua
/usr/lib/lua/vendorextensions/ubusPlugin.lua
/usr/lib/lua/vendorextensions/eventWatcher.lua
/usr/lib/lua/vendorextensions/msg_encode.lua
/usr/lib/lua/vendorextensions/agent_ipc_handler.lua
/etc/init.d/multiap_vendorextensions
