# usr/lib/opkg/info/multiap_controller.list
/usr/bin/multiap_controller
/etc/init.d/multiap_controller
