# usr/lib/opkg/info/multiap_agent.list
/etc/init.d/multiap_agent
/lib/multiap/multiap-vendor-ext-uci-to-procd-env
/usr/bin/multiap_agent
