# usr/lib/opkg/info/multiap_agent.list
/etc/init.d/multiap_agent
/usr/bin/multiap_agent
