# usr/lib/opkg/info/multiap_controller.list
/usr/bin/multiap_controller
/lib/multiap/multiap-controller-uci-to-procd-env
/lib/multiap/multiap-ctrl-ve-uci-to-procd-env
/usr/bin/map_common_reset_credentials.sh
/usr/bin/map_reset_to_default.sh
/etc/init.d/multiap_controller
