# usr/lib/opkg/info/mvfs.list
/usr/bin/realpath
/usr/bin/mvfs
/etc/init.d/mvfs
/usr/bin/readdir
/usr/lib/libmvfsl.so
/usr/bin/readview
