# usr/lib/opkg/info/mwan.list
/etc/init.d/mwan
/etc/init.d/pre-mwan
/etc/config/mwan
/etc/hotplug.d/iface/15-mwan
/usr/lib/mwan/functions.sh
/usr/lib/mwan/setup.sh
