# usr/lib/opkg/info/nfcapp.list
/etc/init.d/nfcd
/usr/bin/nfcd.lua
/usr/sbin/nfcd
