# usr/lib/opkg/info/nginx.list
/etc/nginx/nginx.conf
/usr/sbin/nginx
/etc/init.d/nginx
/etc/nginx/mime.types
/lib/upgrade/keep.d/nginx
