# usr/lib/opkg/info/ngwfdd.list
/usr/share/ngwfdd/mod_process.lua
/usr/share/ngwfdd/mod_voice.lua
/usr/share/ngwfdd/mod_syslog.lua
/usr/share/ngwfdd/mod_dsl.lua
/usr/share/ngwfdd/mod_wan.lua
/usr/share/ngwfdd/lib/gwfd-common.lua
/usr/share/ngwfdd/mod_dev_connected.lua
/usr/share/ngwfdd/mod_meminfo.lua
/usr/share/ngwfdd/mod_stats.lua
/etc/init.d/ngwfdd
/usr/bin/ngwfdd
/usr/share/ngwfdd/mod_video.lua
/etc/config/ngwfdd
/usr/share/ngwfdd/lib/timed_write.sh
/usr/share/ngwfdd/mod_device.lua
/usr/share/ngwfdd/mod_lte.lua
