# usr/lib/opkg/info/ngwfdd.list
/usr/share/ngwfdd/mod_process.lua
/usr/share/ngwfdd/mod_slabinfo.lua
/usr/share/ngwfdd/mod_fhcd.lua
/usr/share/ngwfdd/mod_voice.lua
/usr/share/ngwfdd/mod_syslog.lua
/usr/share/ngwfdd/mod_dsl.lua
/usr/share/ngwfdd/mod_wan.lua
/usr/lib/lua/gwfd/common.lua
/usr/share/ngwfdd/mod_dev_connected.lua
/usr/lib/ngwfdd/sanitize_config
/usr/share/ngwfdd/mod_cpustats.lua
/usr/share/ngwfdd/mod_meminfo.lua
/usr/bin/pgroup
/etc/init.d/ngwfdd
/usr/share/ngwfdd/mod_conntrack.lua
/usr/bin/ngwfdd
/usr/share/ngwfdd/mod_video.lua
/usr/share/ngwfdd/lib/timed_write.sh
/usr/share/ngwfdd/mod_device.lua
/usr/share/ngwfdd/mod_lte.lua
