# usr/lib/opkg/info/nqtools.list
/etc/config/samba
/usr/bin/nqnd
/etc/init.d/mountd
/etc/hotplug.d/mount/00-mounts
/etc/config/mounts
/lib/functions/samba_printer.sh
/usr/bin/nqcsctrl
/etc/hotplug.d/mount/00-samba
/lib/functions/contentsharing.sh
/etc/hotplug.d/usb/21-samba-printer
/usr/bin/nqcs
/etc/init.d/nqe
