# usr/lib/opkg/info/odhcpd.list
/usr/sbin/odhcpd
/etc/init.d/odhcpd
/usr/sbin/odhcpd-update
