# usr/lib/opkg/info/openvpn-openssl.list
/etc/config/openvpn
/etc/init.d/openvpn
/usr/sbin/openvpn
/lib/upgrade/keep.d/openvpn
