# usr/lib/opkg/info/opticald.list
/etc/init.d/opticald
/usr/sbin/opticald
