# usr/lib/opkg/info/pon-caldata.list
/usr/sbin/bosa_caldata.sh
/etc/init.d/pon_caldata
/etc/bosa.dat
