# usr/lib/opkg/info/pppoe-relay-tch.list
/etc/hotplug.d/net/50-pppoe-relay-tch
/usr/lib/pppoe-relay-tch/setup.sh
/etc/init.d/pppoe-relay-tch
