# usr/lib/opkg/info/pptpd.list
/usr/sbin/pptpd
/usr/sbin/pptpctrl
/etc/init.d/pptpd
/etc/ppp/options.pptpd
/etc/pptpd.conf
/usr/lib/pptpd/pptpd-logwtmp.so
/usr/sbin/bcrelay
/etc/config/pptpd
