# usr/lib/opkg/info/qos-tch.list
/etc/init.d/qos
/usr/bin/qos
/etc/hotplug.d/iface/00-qos
/usr/lib/qos/generate.sh
/etc/config/qos
