# usr/lib/opkg/info/qos-tch.list
/etc/init.d/qos
/etc/init.d/dscp2pbit
/usr/lib/iptables/libxt_MARK_REMAP.so
/usr/bin/qos
/etc/hotplug.d/iface/00-qos
/usr/lib/qos/generate.sh
/etc/hotplug.d/iface/10-dscp2pbit
/etc/config/qos
/etc/config/dscp2pbit
