# usr/lib/opkg/info/quantenna.list
/qtn/qtn-linux.lzma
/qtn/qtn_custo.env
/usr/sbin/qtn_call_qcsapi.sh
/qtn/qtn-uboot
/etc/init.d/quantenna
/usr/sbin/qtn_qcsapi_ip
/qtn/qtn_custo.sh
/usr/sbin/qtn_qcsapi_eth
