# usr/lib/opkg/info/quantenna.list
/qtn/qtn-linux.lzma
/usr/sbin/qtn_call_qcsapi.sh
/qtn/qtn-uboot
/etc/init.d/quantenna
/usr/sbin/qtn_qcsapi_ip
/usr/sbin/qtn_qcsapi_eth
