# usr/lib/opkg/info/random-seed.list
/lib/preinit/99_random_seed
/usr/sbin/random_seed
/etc/init.d/random_seed
