# usr/lib/opkg/info/reset-to-factory-default.list
/lib/mount_root/20_rtfd_remove
/usr/bin/rtfd_restore
/usr/bin/rtfd_preserve
/usr/bin/rtfd
/etc/init.d/restore
/etc/init.d/rtfd
/lib/mount_root/00_rtfd_all
/lib/mount_root/80_rtfd_restore
