# usr/lib/opkg/info/reset-to-factory-default.list
/usr/bin/rtfd_restore
/usr/bin/rtfd_preserve
/usr/bin/rtfd
/etc/init.d/restore
