# usr/lib/opkg/info/rp-pppoe-relay.list
/usr/sbin/pppoe-relay
/lib/upgrade/keep.d/rp-pppoe-relay
/etc/init.d/pppoe-relay
/etc/uci-defaults/pppoe-relay.default
