# usr/lib/opkg/info/rp-pppoe-relay.list
/usr/sbin/pppoe-relay
/etc/hotplug.d/iface/60-pppoe-relay
/etc/init.d/pppoe-relay
