# usr/lib/opkg/info/samba3.list
/etc/config/samba
/etc/samba/smbpasswd
/etc/samba/smb.conf.template
/etc/init.d/samba
/usr/sbin/smbpasswd
/usr/sbin/smbd
