# usr/lib/opkg/info/samba-tch.list
/etc/config/samba
/etc/samba/smbpasswd
/etc/samba/smb.conf.template
/etc/config/mounts
/etc/init.d/samba
/lib/functions/samba_printer.sh
/etc/init.d/mountd-procd
/etc/hotplug.d/mount/00-samba
/usr/sbin/smbpasswd
/usr/sbin/smbd
/etc/hotplug.d/usb/21-samba-printer
/etc/hotplug.d/mount/20-mounts
/usr/bin/dfree
