# usr/lib/opkg/info/sfp-get.list
/usr/sbin/sfp_get.sh
/etc/init.d/sfp-hisr-pppoe-graceful-restart
/etc/hotplug.d/platform/00-sfp
