# usr/lib/opkg/info/sfp-i2cctl.list
/etc/config/optical
/usr/sbin/sfpi2cctl
/etc/init.d/sfpmon
/etc/hotplug.d/platform/00-sfp
