# usr/lib/opkg/info/snmpd.list
/etc/init.d/snmpd
/etc/default/snmpd
/etc/config/snmpd
/etc/uci-defaults/91-net-snmp
/lib/functions/firewall-net-snmp.sh
/usr/sbin/snmpd
/etc/snmp/snmpd.conf
