# usr/lib/opkg/info/snmpd.list
/etc/init.d/snmpd
/usr/sbin/snmpd
/etc/snmp/snmpd.conf
