# usr/lib/opkg/info/socat.list
/etc/init.d/socat
/usr/bin/socat
/etc/config/socat
