# usr/lib/opkg/info/sudo.list
/usr/lib/sudo/libsudo_util.so.0.0.0
/usr/lib/sudo/sudo_noexec.so
/etc/sudoers
/usr/lib/sudo/sudoers.so
/usr/bin/sudo
/usr/lib/sudo/group_file.so
/etc/init.d/sudo
/usr/sbin/visudo
/lib/upgrade/keep.d/sudo
/usr/lib/sudo/libsudo_util.so
/usr/lib/sudo/system_group.so
/usr/lib/sudo/libsudo_util.so.0
