# usr/lib/opkg/info/sudo.list
/etc/sudoers
/usr/bin/sudo
/etc/init.d/sudo
/usr/sbin/visudo
/usr/lib/sudo_noexec.so
