# usr/lib/opkg/info/supervision-tch.list
/etc/init.d/supervision
/usr/bin/supervision.lua
