# usr/lib/opkg/info/thermal-protection.list
/etc/init.d/thermalProtection
/usr/bin/thermalProtection
