# usr/lib/opkg/info/tmsctl.list
/usr/bin/tmsctl
/etc/hotplug.d/net/10_ethoam
/etc/init.d/ethoam
