# usr/lib/opkg/info/tod.list
/sbin/time_change_monitor.lua
/etc/init.d/tod
/usr/lib/lua/todfw/todtimer.lua
/usr/lib/lua/todfw/init.lua
/lib/functions/tod.sh
/usr/lib/lua/todfw/todasync.lua
/usr/lib/lua/todfw/todhelper.lua
/usr/lib/lua/todfw/todaction.lua
/sbin/tod.lua
