# usr/lib/opkg/info/trafficmon.list
/sbin/trafficmon.lua
/etc/init.d/trafficmon
