# usr/lib/opkg/info/transformer-tch.list
/usr/lib/lua/transformer/mapper/multiuci.lua
/usr/lib/lua/transformer/eventsource/ubus.lua
/usr/lib/lua/transformer/eventsource/uci.lua
/usr/lib/lua/datamodel.lua
/usr/lib/lua/transformer/mapper/ubus.lua
/usr/lib/lua/transformer/mapper/nwcommon.lua
/etc/init.d/transformer
/usr/lib/lua/transformer/eventhor.lua
/usr/lib/lua/transformer/mapper/ucilist.lua
/usr/lib/lua/transformer/xref.lua
/usr/lib/libtransformer.so
/usr/bin/dmdump
/usr/lib/lua/transformer/persistency/conversion/upgrade_to_1.lua
/usr/lib/lua/transformer/lock.lua
/usr/lib/lua/transformer/msg_decode.lua
/usr/lib/lua/transformer/maphelper.lua
/usr/lib/lua/transformer/persistency.lua
/usr/lib/lua/datamodel-bck.lua
/usr/lib/lua/transformer/persistency/conversion/upgrade_to_2.lua
/usr/lib/lua/transformer/mapper/uci_1to1.lua
/usr/lib/lua/transformer/mapper/ucihelper.lua
/usr/lib/lua/transformer/msg.lua
/usr/lib/lua/transformer/commitapply.lua
/usr/lib/lua/transformer/navigation.lua
/usr/lib/lua/transformer/logger.lua
/usr/lib/lua/transformer/persistency/db.lua
/usr/lib/lua/transformer/typecheck.lua
/usr/lib/lua/lasync.so
/usr/lib/lua/transformer/pathfinder.lua
/usr/lib/lua/transformer/persistency/convert.lua
/usr/lib/lua/transformer/mapload.lua
/usr/bin/transformer
/usr/lib/lua/transformer/mapper/command.lua
/usr/lib/lua/transformer/fault.lua
/usr/lib/lua/transformer/msg_encode.lua
/usr/lib/lua/transformer/alias.lua
/usr/lib/lua/transformer/typestore.lua
/usr/lib/lua/transformer/mapper/multiroot.lua
/usr/lib/lua/lcrypto.so
/usr/lib/lua/transformer/mapper/simpleuci.lua
/usr/lib/lua/transformer/api.lua
