# usr/lib/opkg/info/txmember.list
/etc/init.d/txmember
/usr/bin/txmember
/usr/lib/lua/wifidoctor/ipping.lua
/usr/share/transformer/commitapply/rpc_qeo_txmember.ca
/usr/lib/lua/wifidoctor/acs.lua
/etc/config/wifidoctor_station_importance
