# usr/lib/opkg/info/urlfilterd-tch.list
/etc/config/parental
/bin/urlfilterd-tch
/etc/init.d/urlfilterd-tch
