# usr/lib/opkg/info/wansensing.list
/usr/lib/lua/wansensingfw/timedevent.lua
/usr/lib/lua/wansensingfw/init.lua
/usr/lib/lua/wansensingfw/l3wasmachine.lua
/usr/lib/lua/wansensingfw/repeatedcheck.lua
/etc/init.d/wansensing
/sbin/wansensing.lua
/usr/lib/lua/wansensingfw/scripthelpers.lua
/usr/lib/lua/wansensingfw/l2wasmachine.lua
/usr/lib/lua/wansensingfw/washelper.lua
/usr/lib/lua/wansensingfw/wasworker.lua
/usr/lib/lua/wansensingfw/wasasync.lua
