# usr/lib/opkg/info/watchdog.list
/usr/bin/repair.sh
/usr/sbin/watchdog-tch
/lib/upgrade/watchdog.sh
/etc/init.d/watchdog-tch
