# usr/lib/opkg/info/weburl-nfq.list
/etc/init.d/weburl
/etc/weburl_cfg_tmpl
/usr/sbin/urlfilterd
