# usr/lib/opkg/info/webui-app-certpairing.list
/www/lua/generic/pairing/certificate.lua
/www/lua/generic/pairing/dbcert.lua
/www/lua/generic/pairing/dm_pairing.lua
/usr/share/transformer/mappings/rpc/generic_app.hostbyip.map
/etc/init.d/generic_pnpd
/usr/lib/lua/wireless/generic_app/wpsevent.lua
/www/lua/generic/pairing/chrono.lua
/usr/lib/lua/generic/mappings/hostbyip.lua
/usr/share/transformer/mappings/rpc/generic_app.app.wps.map
/usr/share/transformer/mappings/igd/X_000E50_PushNotifications.map
/usr/share/transformer/mappings/rpc/generic_app.app.pairing.map
/www/lua/generic/app.lua
/usr/bin/generic_pnpd
/usr/share/transformer/mappings/uci/generic_app.map
/usr/share/transformer/mappings/rpc/sys.generic_app.map
/usr/share/transformer/commitapply/wpsevent.ca
/usr/lib/lua/wireless/generic_app/wpsrequest.lua
/www/lua/generic/pairing/registry.lua
/etc/init.d/cert
/usr/sbin/pairing_button.sh
/usr/lib/lua/wireless/generic_app/wpsstate.lua
/www/lua/generic/generic_app/identify.lua
/www/lua/generic/pairing/button.lua
/etc/ssl/data/key.VANT-9
/usr/share/transformer/mappings/rpc/generic_app.map
/etc/nginx/loc_pairing.conf
/usr/lib/lua/wireless/generic_app/sender.lua
/usr/share/transformer/commitapply/uci_generic_pnpd.ca
/www/lua/generic/pairing/pairer.lua
/usr/share/transformer/mappings/rpc/generic_app.app.map
/etc/ssl/data/key.VBNT-6
/etc/nginx/server_webservice.conf
/etc/ssl/data/key.VBNT-W
/usr/share/transformer/mappings/igd/X_GENERIC_APP.map
/etc/uci-defaults/tch_0081-generic_app
/www/lua/generic/generic_app/setcert.lua
/usr/share/transformer/commitapply/rpc_app_pairing.ca
/www/lua/generic/generic_app/auth.lua
