# usr/lib/opkg/info/wfa-testsuite-daemon.list
/usr/sbin/wfa-testsuite-daemon
/usr/sbin/wireless_wfa_test_exec.sh
/data/ofdma.sh
/etc/init.d/wfa-testsuite-daemon
/usr/sbin/wireless_wfa_test_init.sh
