# usr/lib/opkg/info/wfa-testsuite-daemon.list
/usr/sbin/wfa-testsuite-daemon
/etc/init.d/wfa-testsuite-daemon
