# usr/lib/opkg/info/wifi-conductor.list
/usr/bin/wifi_conductor.sh
/usr/bin/wifi-conductor-cli
/lib/wifi-conductor/wifi-conductor-uci-to-procd-env
/etc/init.d/wifi-conductor
/usr/lib/mosquitto/wd_auth_plugin.so
/usr/bin/wifi-conductor
