# usr/lib/opkg/info/wifi-conductor.list
/usr/bin/wifi_conductor.sh
/usr/bin/wifi-conductor-cli
/etc/init.d/wifi-conductor
/usr/bin/wifi-conductor
