# usr/lib/opkg/info/wifi-services-agent.list
/usr/bin/wifi-doctor-agent
/etc/init.d/ewifi
/etc/init.d/wifi-doctor-agent
