# usr/lib/opkg/info/wlan.list
/usr/sbin/wps_virtual_button_pressed.sh
/usr/sbin/acl_button_pressed.sh
/lib/netifd/wireless/lantiq.sh
/usr/sbin/mon_rxframe_rxfifo_ovrflow.sh
/usr/sbin/get_wlan_feature.sh
/usr/sbin/wireless_caldata.sh
/etc/init.d/hostapd
/usr/sbin/wireless_restore_default_config.sh
/lib/wireless/get_remote_network.sh
/lib/wireless/init_qcacld.sh
/usr/sbin/hostapd_env.sh
/usr/sbin/hostapd
/usr/sbin/wifi_on.sh
/usr/sbin/parse_pa_trimming.sh
/etc/wireless_acs.conf
/usr/sbin/wireless_get_mac_addr.sh
/lib/wireless/fix_uci_config.sh
/usr/sbin/wireless_autochannel.sh
/etc/uci-defaults/tch_0099-wireless_force_channel36_on_dev
/usr/sbin/hostapd_keygen
/lib/wireless/init_lantiq.sh
/usr/sbin/wireless_get_radio_type.sh
/usr/sbin/hostapd_cli
/usr/sbin/wifionoff.sh
/usr/sbin/wps_button_pressed.sh
/lib/netifd/wireless/broadcom.sh
/lib/wireless/init_broadcom.sh
/usr/sbin/wireless_odm_test_init.sh
/lib/wireless/init_atheros.sh
/usr/sbin/wireless_txtest.sh
/lib/netifd/wireless/qcacld.sh
/lib/wireless/fatal_error.sh
/usr/sbin/set_wlan_feature.sh
/usr/sbin/wireless_init_uci_env.sh
/usr/sbin/nvram2rip.sh
/usr/sbin/wireless_get_overview.sh
/lib/netifd/wireless/atheros.sh
