# usr/lib/opkg/info/wol-tch.list
/etc/hotplug.d/iface/99-wol
/etc/init.d/wol
/lib/functions/firewall-wol.sh
/etc/uci-defaults/19-wol-tch
