# usr/lib/opkg/info/wol-tch.list
/lib/upgrade/keep.d/wol-tch
/etc/init.d/wol
/usr/bin/wol
