# usr/lib/opkg/info/wpa_supplicant_mgmt.list
/etc/init.d/wpa_supplicant_mgmt
/usr/bin/wpa_supplicant_mgmt
