# usr/lib/opkg/info/wpa_supplicant_mgmt.list
/usr/sbin/wireless_ext_manage_ifaces.sh
/etc/init.d/wpa_supplicant_mgmt
/usr/sbin/wps_synchronize_config.sh
/usr/sbin/wireless_wps_ext_easy_install.sh
/usr/bin/wpa_supplicant_mgmt
