# usr/lib/opkg/info/xinetd.list
/etc/xinetd.conf
/usr/sbin/xinetd
/etc/init.d/xinetd
/etc/xinetd.d/cups
