# usr/lib/opkg/info/xinetd.list
/lib/upgrade/keep.d/xinetd
/etc/xinetd.conf
/usr/sbin/xinetd
/etc/init.d/xinetd
