# usr/lib/opkg/info/xl2tpd.list
/etc/xl2tpd/xl2tpd.conf
/lib/netifd/proto/l2tp.sh
/usr/sbin/xl2tpd
/etc/init.d/xl2tpd
/usr/sbin/xl2tpd-control
/etc/ppp/options.xl2tpd
/etc/xl2tpd/xl2tp-secrets
