# usr/lib/opkg/info/xtmctl.list
/usr/bin/xtmctl
/etc/hotplug.d/dsl/10_trigger_xtm
/etc/init.d/xtm
